# sysinit/bootlogs.conf
description "bootlogs: save kernel messages to /var/log/dmesg"

start on sysinit-up and filesystem and started hostname

pre-start script
	test ! -f /var/log/dmesg || savelog -q -p -c 5 /var/log/dmesg

	if which dmesg >/dev/null 2>&1; then
		dmesg -s 524288 > /var/log/dmesg
	elif [ -c /dev/klog ]; then
		dd if=/dev/klog of=/var/log/dmesg count=1024
	fi

	test -f /var/log/dmesg && chgrp adm /var/log/dmesg || :
end script
