# sysinit/hwclock.conf
# hwclock - adjust system clock and timezone
#
# The hwclock task adjusts the system clock when the hardware clock is
# set to localtime (e.g. when dual-booting with Windows), and also
# ensures that the system timezone is set so that timestamps are written
# to FAT devices.

description	"adjust system clock and timezone"

start on starting mountall
stop on starting halt or starting reboot

env UTC=yes
env BADYEAR=no

task

pre-start script
    [ ! -r /etc/default/rcS ] || . /etc/default/rcS
    [ "$UTC" = "yes" ] && tz="--utc" || tz="--localtime"
    [ "$BADYEAR" != "yes" ] || badyear="--badyear"
    hwclock --systz $tz --noadjfile $badyear || true
end script

post-stop script
    test -e /dev/rtc0 || exit 0
    [ ! -r /etc/default/rcS ] || . /etc/default/rcS
    [ "$UTC" = "yes" ] && tz="--utc" || tz="--localtime"
    [ "$BADYEAR" != "yes" ] || badyear="--badyear"
    exec hwclock --rtc=/dev/rtc0 --systohc $tz --noadjfile $badyear
end script
