# sysinit/motd.conf
description "motd: create /run/motd.dynamic"

start on sysinit-up and started hostname and virtual-filesystems

pre-start exec uname -snrvm > /run/motd.dynamic
