# sysinit/motd.conf
description "motd: create /run/motd.dynamic"

start on sysinit-up and started hostname and virtual-filesystems

pre-start
	if test -d "/etc/update-motd.d"; then
		uname -snrvm >/run/motd.dynamic
		stop
	fi
end script

exec run-parts --lsbsysinit /etc/update-motd.d >/run/motd.dynamic
