# sysinit/procps.conf
description "procps - Configure kernel parameters at boot"

start on starting net

pre-start exec /sbin/sysctl --system -q
