# ssh/ssh-keygen.conf
start on starting ssh

task

script
    [ ! -e /etc/ssh/ssh_host_rsa_key ] && \
            ssh-keygen -f /etc/ssh/ssh_host_rsa_key -N '' -t rsa >/dev/null 2>&1
    [ ! -e /etc/ssh/ssh_host_dsa_key ] && \
            ssh-keygen -f /etc/ssh/ssh_host_dsa_key -N '' -t dsa >/dev/null 2>&1
end script
