# mediaserver/forked-daapd.conf
description "DAAP/DACP (iTunes) server, support for AirPlay and Spotify"
author      "Espen Jürgensen <espenjurgensen@gmail.com>"

start on started avahi-daemon
stop on stopping avahi-daemon

# If the server fails try to respawn it, but not more than once in a minute
respawn
respawn limit 1 60

# Run in foreground - if not we would need to set the expect stanza
exec forked-daapd -f
