# mediaserver/mediatomb.conf
description "mediatomb - uPnP media server"
author "Cameron Norman <camerontnorman@gmail.com>"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

pre-start script
	ROUTE_ADD="/sbin/route add -net 239.0.0.0 netmask 255.0.0.0"

	if test -r /etc/default/mediatomb; then
		. /etc/default/mediatomb
	fi

	test -n "$USER" && USER=root
	test -n "$GROUP" && GROUP="$USER"

	LOGFILE=/var/log/mediatomb.log
	touch $LOGFILE
	chown $USER:$GROUP $LOGFILE

	test -n "$INTERFACE" && {
		OPTIONS="$OPTIONS -e $INTERFACE"
		$ROUTE_ADD $INTERFACE
	}

	initctl set-env ARGS="-c /etc/mediatomb/config.xml -u $USER -g $GROUP -l $LOGFILE $OPTIONS"
end script

respawn
exec /usr/bin/mediatomb $ARGS

post-stop script
	ROUTE_DEL="/sbin/route del -net 239.0.0.0 netmask 255.0.0.0"

	if test -r /etc/default/mediatomb; then
                . /etc/default/mediatomb
        fi

	test -n "$INTERFACE" && $ROUTE_DEL $INTERFACE
end script
