# base/ifupdown/etc/init/network-interface-stall.conf
description "network-interface-stall: wait for network-pre checkpoint"
author "Cameron Nemo <camerontnorman@gmail.com>"

# In order to avoid upstart bug LP: #447654, we cannot have an AND
# statement here (with the ORs).  An "and virtual-filesystems" is desired
# here to make sure that the securityfs is mounted, but since each of the
# ORed services already require virtual-filesystems be mounted, this is safe:
start on (starting network-interface
          or starting network-manager)
stop on started network-pre

# In order to handle the lack of upstart feature LP: #568860, we need to
# run multiple times, for each of the above "starting" service instances, or
# else another one might run while we're running, and not wait for us to
# finish.
instance $JOB${INTERFACE:+/}${INTERFACE:-}

task

script
    # don't stall the loopback iface
    test "$INTERFACE" != "lo" || exit 0
    # test if network-pre job exists
    status network-pre >/dev/null 2>&1 || exit 1
    # wait for the job to reach desired state
    while status network-pre | grep -v -q "start/running"; do
        sleep 1
    done
end script
