# base/ifupdown/etc/init/network-pre.conf
description "network-pre: checkpoint for networking"
author "Cameron Nemo <camerontnorman@gmail.com>"

start on (local-filesystems
          and (started udevtrigger
	       or container))
