# base/tty/container-getty.conf
description "start getty jobs in a container"

start on started login and container CONTAINER=lxc

env MAX_TTY=6

pre-start script
    start -n getty PORT=console || true
    for f in /dev/pts/*; do
        test "$f" = "/dev/pts/ptmx" && continue
	test "${f#/dev/pts/}" -lt "$MAX_TTY" || continue
        start -n getty PORT=${f#/dev/} || true
    done
end script
