# base/tty/container-getty.conf
description "start getty jobs in a container"

start on started login and container CONTAINER=lxc

env ttydir=lxc

pre-start script
    start -n getty PORT=console || true
    for f in /dev/"${ttydir}"/*; do
        start -n getty PORT=${f#/dev/} || true
    done
end script
