# base/tty/container-ttys.conf
description "container-ttys: start getty jobs in a container"

start on container and login-up

pre-start script
	for t in $container_ttys; do
		start -q -n getty PORT=${f#/dev/} || true
	done
end script
