# base/tty/login.conf
start on started sysinit
stop on unmounting-filesystem
