# base/tty/login.conf
description "login: login checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on multiuser-up
stop on multiuser-down

emits login-up
emits login-down

pre-start exec initctl emit login-up
post-stop exec initctl emit login-down
