# base/tty/login.conf
start on started multiuser
stop on stopping multiuser
