# base/tty/std-getty.conf
description "start getty jobs (no container)"

start on started login and not-container

env MAX_TTY=6

pre-start script
    for n in $(seq 1 "$MAX_TTY"); do
        test -e "/dev/tty${n}" || continue
        start -n getty PORT="tty${n}" || true
    done
end script
