# base/tty/std-ttys.conf
description "std-ttys: start standard getty jobs"

start on not-container and login-up

env MIN_TTY=2
env MAX_TTY=6

pre-start script
    for n in $(seq "$MIN_TTY" "$MAX_TTY"); do
        test -e "/dev/tty${n}" || continue
        start -q -n getty PORT="tty${n}" || true
    done
end script
