# net/ifupdown/etc/init/network-interface-stall.conf
description "network-interface-stall: wait for net checkpoint"
author "Cameron Nemo <camerontnorman@gmail.com>"

# we can assume that virtual fs are available with these start events
start on (starting network-interface
          or starting network-manager)
stop on net-up

# In order to handle the lack of upstart feature LP: #568860, we need to
# run multiple times, for each of the above "starting" service instances, or
# else another one might run while we're running, and not wait for us to
# finish.
instance $JOB${INTERFACE:+/}${INTERFACE:-}
export INTERFACE

script
    # don't stall loopback
    test "$INTERFACE" != "lo" || exit 0
    # test if net job exists and is running
    stat="$(status net)" || exit 1
    echo "$stat" | grep -v -q "start/running" || exit 0
    # wait for net-up event
    while sleep 60; do
    	:
    done
end script
