# firewall/ipset.conf
description "ipset: load and unload IP sets"

start on starting firewall
stop on stopped firewall

pre-start script
	test -s /etc/iptables/ipset || { stop; exit 0; }
	exec ipset restore -file /etc/iptables/ipset
end script

post-stop script
	/bin/ipset flush
	/bin/ipset destroy
end script
