# firewall/netfilters-persistent.conf
description "nftables-persistent: netfilter persistent configuration"

start on firewall-up
stop on firewall-down

# This is better ??

#start on firewall-up
#env ACTION=start
#instance $ACTION
#usage "ACTION=start|save|flush"
#task
#script
#	test "$ACTION" = "start" -o "$ACTION" = "save" -o "$ACTION" = "flush"
#	test -d /usr/share/netfilter-persistent/plugins.d || exit 0
#	exec run-parts -v -a "$ACTION" /usr/share/netfilter-persistent/plugins.d
#end script

#start on firewall-down
#manual
#task
#exec start nftables-persistent ACTION=flush

pre-start exec /usr/sbin/netfilter-persistent start
post-stop exec /usr/sbin/netfilter-persistent stop
