# firewall/pyroman.conf
description "pyroman - firewall"

start on firewall-up
stop on firewall-down

pre-start script
    test -x /usr/sbin/pyroman || { stop; exit 0; }
    /usr/sbin/pyroman --init
end script
