# cinder/cinder-volume.conf
description "OpenStack Cinder Volume"
author "Thomas Goirand <zigo@debian.org>"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /var/run

pre-start script
	for i in lock run log lib ; do
		mkdir -p /var/$i/cinder
		chown cinder /var/$i/cinder
	done
end script

script
	[ -x "/usr/bin/cinder-volume" ] || exit 0
	DAEMON_ARGS=" --config-file=/etc/cinder/cinder.conf"
	[ -r /etc/default/openstack ] && . /etc/default/openstack
	[ -r /etc/default/$UPSTART_JOB ] && . /etc/default/$UPSTART_JOB
	[ "x$USE_SYSLOG" = "xyes" ] && DAEMON_ARGS="$DAEMON_ARGS --use-syslog"
	[ "x$USE_LOGFILE" != "xno" ] && DAEMON_ARGS="$DAEMON_ARGS --log-file=/var/log/cinder/cinder-volume.log"

	exec start-stop-daemon --start --chdir /var/lib/cinder \
		--chuid cinder:cinder --make-pidfile --pidfile /var/run/cinder/cinder-volume.pid \
		--exec /usr/bin/cinder-volume -- --config-file=/etc/cinder/cinder.conf ${DAEMON_ARGS}
end script

