# net/torrent/transmission-daemon.conf
description "Bittorrent downloader"

start on filesystem and static-network-up
# Transmission needs to notify peers before it dips
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

# Play with this for a lower cpu priority
#nice 10

setuid debian-transmission
setgid debian-transmission

kill timeout 30
exec /usr/bin/transmission-daemon -f --config-dir /var/lib/transmission-daemon/info
