# 2ping.conf
description "2ping - ping listener daemon"

start on filesystem and net-device-up IFACE=lo
stop on runlevel [016]

respawn
setuid nobody
setgid nogroup

pre-start script
    test -x /usr/bin/2ping || { stop; exit 0; }
end script

exec /usr/bin/2ping --listen -q
