# 4store-avahi-restart.conf
description "Restarts 4store daemons when Avahi is started"

start on started avahi-daemon and started 4store-backend

task

exec /sbin/restart 4store-backend
