# 4store-backend.conf
description "4store-backend - RDF database backend"
author "Cameron Norman <camerontnorman@gmail.com>"

# TODO: Avahi?
start on start-foursquare
stop on stopping 4store

expect daemon
respawn

setuid fourstore
setgid fourstore

env KB=default

script
    test ! -s /etc/default/4store || . /etc/default/4store
    exec /usr/bin/4s-backend $KB
end script
