# 4store-httpd.conf
description "4store-httpd - RDF database"
author "Cameron Norman <camerontnorman@gmail.com>"

start on start-foursquare
stop on stopping 4store

expect daemon
respawn

setuid fourstore
setgid fourstore

env KB=default
env PORT=9000

script
    test ! -s /etc/default/4store || . /etc/default/4store
    exec /usr/bin/4s-httpd -p $PORT $KB
end script
