# acpi-support.conf
description "acpi-support - start power management scripts"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

pre-start script
    if [ -f /usr/share/acpi-support/power-funcs ]; then
        . /usr/share/acpi-support/power-funcs
        on_ac_power || /etc/acpi/power.sh true
    fi
    # Source everything in /etc/acpi/start.d
    if [ -d /etc/acpi/start.d ]; then
        for script in /etc/acpi/start.d/*.sh; do
            . "$script"
        done
    fi
end script

post-stop script
    if [ -f /usr/share/acpi-support/power-funcs ]; then
        . /usr/share/acpi-support/power-funcs
        on_ac_power || /etc/acpi/power.sh true
    fi
end script
