# adjtimex.conf
description "adjtimex - set the kernel time variables"
author "Cameron Norman <camerontnorman@gmail.com>"

start on virtual-filesystems

pre-start script
    test -x /sbin/adjtimex || { stop; exit 0; }
    test ! -s /etc/default/adjtimex || . /etc/adjtimex
    /sbin/adjtimex -tick ${TICK:-10000} -frequency ${FREQ:-0}
end script
