# ahcpd.conf
description "ahcpd - Ad-Hoc Configuration Protocol daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/sbin/ahcpd || { stop; exit 0; }

    test ! -s /etc/default/ahcpd || . /etc/default/ahcpd

    test -n $INTERFACES || {
        echo "No interfaces to operate on, not starting."
        stop; exit 0
    }
end script

script
    test ! -s /etc/default/ahcpd || . /etc/default/ahcpd
    exec /usr/sbin/ahcpd -D $DAEMON_ARGS $INTERFACES
end script
