# alsa-utils.conf
description "alsa-utils - compatibility job for sysvinit"

start on startup
stop on runlevel [016] or unmounting-filesystem
