# alsa-state.conf
# alsa-state - Manage sound card state (state-daemon.conf scheme)
#
# One of two methods of saving ALSA sound card states. Uses the
# file /etc/alsa/state-daemon.conf for (re)storing state.

description "Manage sound card state (state-daemon.conf scheme)"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and sound-device-added
stop on runlevel [016]

expect fork
respawn

nice 19

# alsactl expects this env variable, give it to them
env HOME=/var/run/alsa

pre-start script
    [ -f /etc/alsa/state-daemon.conf ] || { stop ; exit 0; }
    [ -d $HOME ] || mkdir -p $HOME
end script

exec alsactl --syslog --background --sched-idle rdaemon

pre-stop exec alsactl --syslog kill save_and_quit
