# and.conf
description "and - AutoNice daemon"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start script
    test -x /usr/sbin/and || { stop; exit 0; }
end script

exec /usr/sbin/and
