# anon-proxy.conf
description "anonimizing proxy"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

env PORT=4001
env CASCADE=mix.inf.tu-dresden.de:443

setuid anon-proxy
setgid nogroup

console none

respawn

script
	if test -r /etc/default/anon-proxy ; then
		. /etc/default/anon-proxy
		test -z "$OPTIONS" || echo "Ignoring OPTIONS variable"
	fi

	exec /usr/sbin/mix -a -j -n $CASCADE -p $PORT
end script
