# apache2.conf
description "apache2 - HTTP web server"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork

# Instancing with a default empty value
env AP_INSTANCE=
export AP_INSTANCE
instance $AP_INSTANCE
usage "AP_INSTANCE - name of Apache instance, see /usr/share/doc/apache2/README.multiple-instances"

pre-start script
	# The suffix will be -$instance if this is not the main server
	suffix="${AP_INSTANCE:+-$AP_INSTANCE}"

	APACHE_CONFDIR=/etc/apache2$suffix
	APACHE_ENVVARS=$APACHE_CONFDIR/envvars
	[ -r $APACHE_ENVVARS ] && . $APACHE_ENVVARS

	${APACHE_ULIMIT_MAX_FILES:-ulimit -n 8192} || echo "Setting ulimit failed. See README.Debian for more information."

	# Create the run directory if necessary
	[ -d ${APACHE_RUN_DIR:-/var/run/apache2} ] || mkdir -p ${APACHE_RUN_DIR:-/var/run/apache2}
 
	# Attempt to create the lock directory
	dir=${APACHE_LOCK_DIR:-/var/lock/apache2}
	owner=${APACHE_RUN_USER:-www-data}
	group=${APACHE_RUN_GROUP:-www-data}
	if [ -e $dir ] && [ ! -d $dir ]; then
		echo "Apache lock path (${dir}) exists, but is not a directory." 
		exit 1
	elif [ -d $dir ] && [ "$(stat -c %U $dir)" != "$owner" ]; then
		echo "Apache lock directory (${dir}) exists, but not owned by ${owner}."
		echo "Please fix manually."
		exit 1
	else
		install -m 755 -o $owner -g $group -d $dir
	fi

	# ssl_scache shouldn't be here if we're just starting up.
	# (this is bad if there are several apache2 instances running)
	rm -f ${APACHE_RUN_DIR:-/var/run/apache2}/*ssl_scache*
end script

script
	# The suffix will be -$instance if this is not the main server
	suffix="${AP_INSTANCE:+-$AP_INSTANCE}"
 
	APACHE_CONFDIR=/etc/apache2$suffix
	APACHE_ENVVARS=$APACHE_CONFDIR/envvars
	[ -r $APACHE_ENVVARS ] && . $APACHE_ENVVARS

	APACHE_ARGUMENTS="-d $APACHE_CONFDIR $APACHE_ARGUMENTS"

	exec /usr/sbin/apache2 ${APACHE_ARGUMENTS} -k start
end script
