# apparmor.conf
description "apparmor - compatibility job for sysv"

start on filesystem
stop on stopped rc RUNLEVEL=[06]
