# apparmor-load.conf
description "apparmor-load - Load profiles in /etc/apparmor.d"
author "Cameron Norman <camerontnorman@gmail.com>"

start on (file FILE=/etc/apparmor.d/* EVENT=create
          or file FILE=/etc/apparmor.d/* EVENT=modify)

instance $MATCH
usage "MATCH - the full path of the AppArmor profile to load"

task

apparmor load $MATCH
