# apparmor-teardown.conf
description "apparmor-teardown - Unload running AppArmor profiles"

start on stopped rc RUNLEVEL=[06]

task

script
    [ -f /lib/apparmor/functions ] && . /lib/apparmor/functions \
    || { stop; exit 0; }
    running_profile_names | while read profile; do
        unload_profile "$profile" || continue
    done
end script
