# arangodb.conf
description "ArangoDB - the multi-purpose NoSQL database"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

respawn

env USER=arangodb
env GROUP=arangodb
env CONF=/etc/arangodb/arangod.conf

pre-start script
	install -o $USER -g $GROUP -d /var/run/arangodb

	dirs="/var/log/arangodb /var/lib/arangodb /var/lib/arangodb-apps"
	for dir in $dirs; do
		test -d $dir && chown -R $USER:$GROUP $dir
	done

	/usr/sbin/arangod -c $CONF --uid $USER --gid $GROUP \
	                  --check-version || {
		echo "Database version check failed."
		echo "Maybe need to run arangod with the upgrade option?"
		exit 2
	}
end script

exec /usr/sbin/arangod -c $CONF --uid $USER --gid $GROUP \
                       --temp-path "/var/tmp/arangod"
