# atftpd.conf
description "Trivial File Transfer Protocol Server."

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

pre-start script
	USE_INETD=true
	if test -r /etc/default/atftpd; then
		. /etc/default/atftpd
		:
	fi
	test "$USE_INETD" != "true" || { stop; exit 0; }
end script

script
	if test -r /etc/default/atftpd; then
		. /etc/default/atftpd
		:
	fi

	exec /usr/sbin/atftpd --no-fork $OPTIONS
end script
