# avahi-dnsconf.conf
description "avahi-dnsconfd - Avahi Unicast DNS Configuration Daemon"

start on remote-filesystems and started avahi-daemon
stop on (unmounted-remote-filesystems
         or stopping avahi-daemon
         or runlevel [016])

expect daemon
respawn

pre-start script
    if [ -s /etc/localtime ]; then
        [ -d /etc/avahi/etc ] || mkdir -p /etc/avahi/etc
        cp -fp /etc/loacaltime /etc/avahi/etc
    fi
end script

exec avahi-dnsconfd -D
