# babeld.conf
description "ad-hoc network routing daemon"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

env INTERFACES=""
env DAEMON_ARGS="-S /var/lib/babeld/state"

pre-start script
	if test -r /etc/default/babeld ; then
		. /etc/default/babeld
		:
	fi

	# Exit if there are no interfaces
	if [ -z "$INTERFACES" ]; then
		echo "No interfaces to operate on"
		stop; exit 0
	fi
end script

script
	if test -r /etc/default/babeld ; then
		. /etc/default/babeld
		:
	fi

	exec /usr/sbin/babeld $DAEMON_ARGS $INTERFACES
end script
