# beanstalkd.conf
description "beanstalkd - simple, fast work queue"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

setuid beanstalkd
setgid beanstalkd

respawn
script
	if test -r /etc/default/beanstalkd; then
		. /etc/default/beanstalkd
	fi
	exec /usr/bin/beanstalkd -l ${BEANSTALKD_LISTEN_ADDR:-127.0.0.1} -p ${BEANSTALKD_LISTEN_PORT:-11300} $BEANSTALKD_EXTRA
end script
