# bind9.conf
description "BIND Domain Name Server"

start on (remote-filesystems
          and started networking
          and started rsyslog)
stop on (runlevel [016]
         or stopping networking
         or stopping rsyslog)

exec named -f -u bind

pre-stop exec rndc stop
