# bind9.conf
description "BIND Domain Name Server"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

apparmor load /etc/apparmor.d/usr.sbin.named

exec /usr/sbin/named -f -u bind

post-start script
    [ -x /sbin/resolvconf ] || exit 0
    echo nameserver 127.0.0.1 | /sbin/resolvconf -a lo.named
end script

pre-stop exec rndc stop

post-stop script
    [ -x /sbin/resolvconf ] || exit 0
    /sbin/resolvconf -d lo.named
end script
