# bind9-resolvconf.conf
description "Use resolvconf to achieve local BIND"

start on started bind9
stop on stopping bind9

pre-start script
    [ -x /sbin/resolvconf ] || { stop; exit 0; }
    echo nameserver 127.0.0.1 | /sbin/resolvconf -a lo.named
end script

pre-stop exec /sbin/resolvconf -d lo.named
