# bind9.conf
description "BIND Domain Name Server"

start on filesystem
stop on runlevel [016]

exec named -f -u bind

pre-stop exec rndc stop
