# binkd.conf
description "binkd - transfer files between two Fidonet systems over TCP/IP"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

respawn
setuid ftn
setgid ftn

script
	grep -sq '^binkp[[:space:]]' /etc/inetd.conf && {
		echo "Not started, started by inetd"; exit 0
	}

	exec /usr/sbin/binkd -q /etc/binkd/binkd.cfg
end script
