# bitz-server.conf
description "bitz ICAP server"
author "Cameron Norman <camerontnorman@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

exec /usr/sbin/bitz-server /etc/bitz/bitz-server.conf
