# bluetooth.conf
description "bluez - Bluetooth subystem starting and stopping"

start on filesystem and started dbus
stop on stopping dbus

expect fork
respawn

exec bluetoothd
