# bluetooth.conf
description "bluez - bluetooth device daemon"

start on started dbus
stop on stopping dbus

# For bluez 4
#expect fork

respawn
exec /usr/sbin/bluetoothd
