# bootlogs.conf
description "bootlogs - handle log files at boot up"

start on filesystem
stop on runlevel [06]

pre-start script
    # Save kernel messages in /var/log/dmesg
    if [ -f /var/log/dmesg ]; then
        savelog -q -p -c 5 /var/log/dmesg
        chgrp adm /var/log/dmesg || :
    fi
    
    if which dmesg >/dev/null 2>&1; then
        dmesg -s 524288 > /var/log/dmesg 
    elif [ -c /dev/klog ]; then
        dd if=/dev/klog of=/var/log/dmesg
    fi
end script
