# bootparamd.conf
description "bootparamd - boot parameter server"

start on filesystem
stop on runlevel [016]

expect fork
exec /usr/sbin/rpc.bootparamd -s
