# bootchart.conf
description "Boot Chart daemon"

start on startup
stop on (stopped rc RUNLEVEL=[2345]
         or bootchart-timeout
         or runlevel [016])

exec /lib/bootchart/bootchart-collector -r 50

pre-stop exec /sbin/bootchartd stop
