# bootmail.conf
# bootmail - Send signed email on startup and shutdown

description "Bootmail"
author "Dustin Kirkland <dustin.kirkland@gmail.com>"

start on stopped rc RUNLEVEL=[2345]
stop on starting rc RUNLEVEL=[06]

pre-start script
	if [ ! -e /var/run/.bootmail-just-installed ]; then
		# Ensure that we run this with a root environment
		su  -c "/usr/sbin/bootmail" -l root
	fi
end script

pre-stop script
	# Ensure that we run this with a root environment
	su -c "/usr/sbin/bootmail --shutdown" -l root
end script

# Fun number; both prime and fibonacci
exec sleep 514229d
