# bootchart-done.conf
description "timeout for bootchartd"

start on runlevel [2345]
stop on stopped bootchart

exec sleep 45

post-start exec /sbin/initctl emit -n bootchart-timeout
