# brltty.conf
description "brltty - Braille terminal driver"
author "Cameron Norman <camerontnorman@gmail.com>"

start on virtual-filesystems

respawn
reload signal 0
nice -10
oom score -900

pre-start script
    test -x /sbin/brltty || { stop; exit 0; }
end script

exec /sbin/brltty -q -n
