# canna.conf
description "canna - Japanese input server"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start script
    test -x /usr/sbin/cannaserver || { stop; exit 0; }
end script

script
    OPTS="-u canna"

    if test -r /etc/default/canna; then
        . /etc/default/canna
    fi
    
    if test "$CANNA_ENABLE_INET" = "yes"; then
        OPTS="$OPTS -inet"
    fi
    
    exec /usr/sbin/cannaserver $OPTS
end script

pre-stop exec /usr/bin/cannakill
