# certmonger.conf
# certmonger
#
# certmonger is a D-Bus-based service which attempts to simplify 
# interaction with certifying authorities (CAs) on networks which use 
# public-key infrastructure (PKI).

description	"Certmonger"

start on net-device-up
stop on runlevel [06]

expect daemon
respawn

exec certmonger
