# clamav-milter.conf
description "ClamAV-milter - sendmail and Postfix virus scanning"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

env CONF=/etc/clamav/clamav-milter.conf

pre-start script
	test -x /usr/sbin/clamav-milter || { stop; exit 0; }

	test -f "$CONF" || {
		echo "There is no configuration file for clamav-milter."
		echo "Please either dpkg-reconfigure clamav-milter, or copy the example from"
		echo "/usr/share/doc/clamav-milter/examples/ to $CONF."
		exit 6
	}

	get_config_key() {
                awk -v wanted_key="$1" '$1 == wanted_key { print $2 }' "$CONF" | head -n1
        }

	# Get variables from config necessary for prep work
	USER="$(get_config_key User)"
	PIDFILE="$(get_config_key PidFile)"

	# Defaults if they were not set
	test -n "$USER" || USER=clamav
	test -n "$PIDFILE" || PIDFILE=/run/clamav/clamav-milter.pid

	initctl set-env USER="$USER"

	# We cannot have the directory of the pidfile be owned by the milter's 
	# user, because clamd may need it and may run as a different user.
	# So, just chown the pidfile specifically.
	mkdir -p $(dirname "$PIDFILE")
	touch "$PIDFILE"
	chown "$USER" "$PIDFILE"
end script

expect fork
respawn
exec start-stop-daemon -S -q -c "$USER" -x /usr/sbin/clamav-milter -- -c "$CONF"
