# clamav-freshclam.conf
description "ClamAV-freshclam - update ClamAV virus database"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

env CONF=/etc/clamav/freshclam.conf

pre-start script
	test -x /usr/bin/freshclam || { stop; exit 0; }

	test -f "$CONF" || {
		echo "There is no configuration file for freshclam."
		echo "Please either dpkg-reconfigure clamav-freshclam, or copy the example from"
		echo "/usr/share/doc/clamav-freshclam/examples/ to $CONF."
		exit 6
	}

	get_config_key() {
                awk -v wanted_key="$1" '$1 == wanted_key { print $2 }' "$CONF" | head -n1
        }

	# Get variables from config necessary for prep work
	DBOWNER="$(get_config_key DatabaseOwner)"
	DBDIR="$(get_config_key DatabaseDirectory)"
	PIDFILE="$(get_config_key PidFile)"

	# Defaults if they were not set
	test -n "$DBOWNER" || DBOWNER=clamav
	test -n "$DBDIR" || DBDIR=/var/lib/clamav
	test -n "$PIDFILE" || PIDFILE=/run/clamav/clamd.pid

	initctl set-env DBOWNER="$DBOWNER"

	# We cannot have the directory of the pidfile be owned by freshclam's 
	# user, because clamd may need it and may run as a different user.
	# So, just chown the pidfile specifically.
	test -f $PIDFILE || {
		mkdir -p $(dirname "$PIDFILE")
		touch "$PIDFILE"
		chown "$DBOWNER" "$PIDFILE"
	}

	# Do not set group; it is derived from clamd's user
	install -o "$DBOWNER" -m 750 -d "$DBDIR"

	if test -e /etc/cron.*/freshclam; then
		stop; exit 0
	fi
end script

expect fork
respawn
exec start-stop-daemon -S -q -c "$DBOWNER" -x /usr/bin/freshclam -- -d --config-file="$CONF"
