# cobbler.conf
# cobbler - provisioning service
#
# Cobbler is a provisioning service

description	"Cobbler"
author "Dustin Kirkland <kirkland@canonical.com>"

start on filesystem and net-device-up
stop on runlevel [016]

respawn

# To add options to your daemon, edit the line below:
exec /usr/bin/cobblerd --no-daemonize
