# configure-printer.conf
description "configure-printer - use udev to configure printers that are added"

start on usb-device-added ID_MODEL=USB_Printer

instance $SUBSYSTEM-$BUSNUM-$DEVNUM

task

exec /lib/udev/udev-configure-printer add $UPSTART_INSTANCE
