# conmux.conf
description "CONMUX registry"

start on runlevel [2345]
stop on runlevel [!2345]

script
    touch /var/run/conmux-registry
    exec /usr/sbin/conmux-registry 63000 /var/run/conmux-registry \
        >>/var/log/conmux/registry.log 2>&1
end script

post-start script
    for f in /etc/conmux/*.cf; do
        start conmux-daemon "CONMUX_CONFIG=$f" &
    done
end script
