# console-setup.conf
description "SysV compat because the bridge does not do well with tasks"

start on started console-setup-set
