# console-setup-set.conf
description "Set up the console font and keyboard"

start on (virtual-filesystems
          or starting rcS
          or starting mountall-shell
          or starting plymouth-splash)

console owner

task

exec setupcon
