# console-setup-save.conf
description "Save the keymap and fonts"

start on started console-setup and filesystem

task

exec setupcon --save-only
