# corosync-notifyd.conf
description "corosync-notifyd - Log Corosync event notifications"
author "Cameron Norman <camerontnorman@gmail.com>"

start on started corosync
stop on stopping corosync

expect fork
respawn

pre-start script
    test -x /usr/sbin/corosync-notifyd || { stop; exit 0; }
end script

exec /usr/sbin/corosync-notifyd -l
post-start exec touch /var/lock/corosync-notifyd
