# corosync.conf
description "corosync - Corosync Cluster Engine"

start on filesystem and net-device-up INTERFACE!=lo
stop on runlevel [016] or deconfiguring-networking

expect fork
respawn

pre-start script
    test -x /usr/sbin/corosync || { stop; exit 0; }
    test -d /var/run/resource-agents || mkdir -p /var/run/resource-agents
end script

exec /usr/sbin/corosync
