# courier-authdaemon.conf
description "courier-authdaemon - authentication for courier daemons"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

respawn
script
	install -o daemon -g daemon -m 775 -d /run/courier
	install -o daemon -g daemon -m 750 -d /run/courier/authdaemon
	. /etc/courier/authdaemonrc
	exec /usr/lib/courier/courier-authlib/authdaemond
end script
