# cubemap.conf
description "cubemap - video reflector for use with VLC"

start on filesystem
stop on runlevel [016]

respawn

pre-start script
    test -x /usr/bin/cubemap || { stop; exit 0; }
    /usr/bin/cubemap --test /etc/cubemap.config
end script

exec /usr/bin/cubemap /etc/cubemap.config
