# cvsd.conf
description "cvsd - run cvs pserver in a secure environment"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start script
    test -x /usr/sbin/cvsd || { stop; exit 0; }
end script

exec /usr/sbin/cvsd -f /etc/cvsd/cvsd.conf
