# cyrus-imapd.conf
description "Cyrus - distributed mail/news daemons"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

env LISTENQUEUE=32
env CONF=/etc/imapd.conf
env MASTERCONF=/etc/cyrus.conf

pre-start script
    test -x /usr/sbin/cyrmaster -a -x /usr/sbin/cyrus || { stop; exit 0; }
    /usr/sbin/cyrus init-helper start
end script

script
    test ! -s /etc/default/cyrus-imapd || . /etc/default/cyrus-imapd
    exec cyrmaster -d -l $LISTENQUEUE -C $CONF -M $MASTERCONF $OPTIONS
end script

post-stop exec /usr/sbin/cyrus init-helper stop
