# dcbd.conf
# dbcd - dbcd daemon
#


description		"dcbd daemon"
author			"Chuck Short <zulcss@ubuntu.com>"

start on runelvel [2345]
stop on runlevel [06]

expect fork
respawn

pre-start script
	modprobe -q dcbnl || true
end script

script
	exec /usr/sbin/dcbd -d $OPTIONS
end script
