# deluged.conf
description "deluged - Deluge BitTorrent Daemon"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

setuid debian-deluged
umask 027

respawn
exec /usr/bin/deluged -d -c /var/lib/deluged/config -l /var/log/deluged/daemon.log -L info
