# dhcpcd.conf
description "dhcpcd - DHCP client daemon"

start on local-filesystems
stop on runlevel [016]

expect fork

pre-start script
    [ -x /sbin/dhcpcd ] || { stop; exit 0; }

    # Check for interface specific instances
    for x in /var/run/dhcpcd-*.pid; do
        [ -f $x ] && { stop; exit 0; }
    done

    if grep -q "^[[:space:]]*iface[[:space:]]*.*[[:space:]]*inet[[:space:]]*dhcp" \
    /etc/network/interfaces; then
        stop; exit 0;
    fi
end script

exec /sbin/dhcpcd -b

pre-stop exec /sbin/dhcpcd -x
