# dictd.conf
description "dictd - dictionary server daemon"

start on filesystem
stop on runlevel [016]

expect fork
respawn

# Default that will be overridden by the default file (he)
env RUN_MODE=daemon

pre-start script
    [ -x /usr/sbin/dictd ] || { stop; exit 0; }
    [ -s /etc/default/dictd ] && . /etc/default/dictd
    [ "$RUN_MODE" = "daemon" ] || { stop; exit 0; }
end script

script
    [ -s /etc/default/dictd ] && . /etc/default/dictd
    exec /usr/sbin/dictd ${DICTD_ARGS}
end script
