# dirmngr.conf
description "dirmngr - download and manage certificate revocation lists"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

expect fork
respawn

pre-start script
    test -x /usr/bin/dirmngr || { stop; exit 0; }
    exec install -o dirmngr -g dirmngr -d /run/dirmngr
end script

exec start-stop-daemon -S -q -c dirmngr -x /usr/bin/dirmngr -- --daemon --sh

post-start script
    test ! -s /etc/default/dirnmngr || . /etc/default/dirmngr
    chmod ${DIRMNGR_SOCKET_MODE:-0777} /run/dirmngr/socket
end script
