# dnsmasq.conf
description "dnsmasq - DNS forwarder and DHCP server"

start on filesystem and started networking
stop on runlevel [016]

apparmor load /etc/apparmor.d/usr.sbin.dnsmasq

env SCRIPT=/etc/init.d/dnsmasq

exec $SCRIPT systemd-exec

post-start exec $SCRIPT systemd-start-resolvconf

post-stop exec $SCRIPT systemd-stop-resolvconf
