# doodled.conf
description "doodled - desktop search engine daemon."
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

respawn
reload signal 0

# Default to be overridden by default file
env DAEMON_OPTS="-D -d /var/lib/doodle/doodle.db -L /var/log/doodle.log /home"

pre-start script
    test -x /usr/bin/doodled || { stop; exit 0; }
end script

script
    if test -r /etc/default/doodled; then
        . /etc/default/doodled
    fi

    exec /usr/bin/doodled $DAEMON_OPTS
end script
