# downtimed.conf
# downtimed - system downtime monitoring daemon
#
# This file should be installed as /etc/init/downtimed.conf
# on architectures which use upstart(8) to bring up system
# daemons.

description	"system downtime monitoring daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork

exec downtimed
