# drac.conf
description "drac - dynamic relay authorization control"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and started rpcbind
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

env OPTIONS="-i -e 30 /var/lib/drac/dracd.db"

pre-start script
    test -x /usr/sbin/rpc.dracd || { stop; exit 0; }
end script

exec /usr/sbin/rpc.dracd $OPTIONS
