# dropbear.conf
description "dropbear - lightweight SSH server"

start on filesystem
stop on runlevel [016]

expect fork
respawn

# Will be overwritten by default file, do not modify here.
env DROPBEAR_PORT=22
env DROPBEAR_RECEIVE_WINDOW=65536
env DROPBEAR_RSAKEY="/etc/dropbear/dropbear_rsa_host_key"
env DROPBEAR_DSSKEY="/etc/dropbear/dropbear_dss_host_key"

pre-start script
    test -x /usr/sbin/dropbear || { stop; exit 0; }
    if test -x /usr/sbin/update-service && update-service --check dropbear
    then
        echo "Controlled with runit, use sv(8) to start"; stop; exit 0
    fi
end script

script
    test ! -s /etc/default/dropbear || . /etc/default/dropbear

    exec /usr/sbin/dropbear -d "$DROPBEAR_DSSKEY" \
                            -r "$DROPBEAR_RSAKEY" \
                            -p "$DROPBEAR_PORT"   \
                            -W "$DROPBEAR_RECEIVE_WINDOW" \
                            ${DROPBEAR_BANNER:+ -b $DROPBEAR_BANNER} \
                            $DROPBEAR_EXTRA_ARGS
end script
