# ebnetd.conf
description "ebnetd - EBNET server"

start on filesystem
stop on runlevel [016]

# It is started from inetd, we are just prepping it
pre-start exec install -o ebnetd -g nogroup -d /var/run/ebnetd
