# edubuntu-server-manager.conf
# edubuntu-server-manager - Management web interface
#
# This job is started at boot time and starts the management
# web interface for Edubuntu Server.

description     "Edubuntu Server"
author          "Stéphane Graber <stgraber@ubuntu.com>"

start on static-network-up
stop on runlevel [!2345]

respawn

env PYTHONDONTWRITEBYTECODE=1

pre-start script
    [ -f /etc/edubuntu-server/manager.cfg ] || stop
end script

exec edubuntu-server-manager
