# ekeyd.conf
description "ekeyd - Entropy Key Daemon"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect daemon
respawn

pre-start script
    test -x /usr/sbin/ekeyd || { stop; exit 0; }
    if test ! -s /etc/entropykey/entropykey.conf; then
        echo "No config, not starting."; stop; exit 0
    fi
end script

exec /usr/sbin/ekeyd
