# exim4-queue-runner.conf
description "exim-queue-runner - separate exim4 queue runner daemon"

start on start-exim
stop on stop-exim

# TODO: readiness
respawn

# Backups in case the default file is unavailable
env QUEUERUNNER="combined"
env QUEUEINTERVAL="30m"

pre-start script
    [ -s /etc/default/exim4 ] && . /etc/default/exim4 || :

    case "$QUEUERUNNER" in
      combined|no|nodaemon|ppp)
        stop; exit 0
        ;;
    esac
end script

script
    [ -s /etc/default/exim4 ] && . /etc/default/exim4 || :

    exec /usr/sbin/exim4 -bdf "-q${QFLAGS}${QUEUEINTERVAL}" \
                         ${COMMONOPTIONS} \
                         ${QUEUERUNNEROPTIONS}
end script
