# exim4-listener.conf
description "exim4-listener - exim4 SMTP port listener"

start on start-exim
stop on stop-exim

# TODO: readiness
respawn

# Backups in case the default file is unavailable
env QUEUERUNNER="combined"
env QUEUEINTERVAL="30m"

pre-start script
    [ -s /etc/default/exim4 ] && . /etc/default/exim4 || :
    [ "$QUEUERUNNER" != "queueonly" ] || { stop; exit 0; }
end script

script
    [ -s /etc/default/exim4 ] && . /etc/default/exim4 || :

    OPTIONS="-bdf ${COMMONOPTIONS} ${SMTPLISTENEROPTIONS}"

    if [ "$QUEUERUNNER" = "combined" ] ; then
        OPTIONS="$OPTIONS -q${QFLAGS}${QUEUEINTERVAL} ${QUEUERUNNEROPTIONS}"
    fi

    exec /usr/sbin/exim4 $OPTIONS
end script
