# farpd.conf
description "farpd - fake ARP daemon"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

# Set NETWORK to a valid address, or set INTERFACE and zero out NETWORK
env NETWORK=unconfigured
env INTERFACE=

pre-start script
    test -x /usr/sbin/farpd || { stop; exit 0; }
    test ! -s /etc/default/farpd || . /etc/default/farpd
    if test "$NETWORK" = "unconfigured" ; then
        echo "farpd has not been configured, do so in the default file."
        exit 1
    fi
end script

script
    test ! -s /etc/default/farpd || . /etc/default/farpd
    exec /usr/sbin/farpd ${INTERFACE:+ -i $INTERFACE} $NETWORK
end script
