# fprobe.conf
description "fprobe - NetFlow collector"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

script
	if test -r /etc/default/fprobe; then
		. /etc/default/fprobe
	fi
	exec /usr/sbin/fprobe ${INTERFACE:+-i$INTERFACE} $OTHER_ARGS ${FLOW_COLLECTOR}
end script
