# freeradius.conf
description "FreeRADIUS - Authentication, Authorization and Accounting server"

start on filesystem and static-network-up
stop on runlevel [016]

expect fork
respawn

env RUNDIR=/var/run/freeradius

pre-start script
    [ -x /usr/sbin/freeradius ] || { stop; exit 0; }
    [ -d $RUNDIR ] || install -d -o freerad -g freerad $RUNDIR
end script

exec /usr/sbin/freeradius

post-stop exec rm -f $RUNDIR/freeradius.pid
