# fso-frameworkd.conf
description "fso-frameworkd - freesmartphone.org framework daemon"

start on started dbus
stop on stopping dbus

respawn

pre-start script
    test -x /usr/bin/fso-frameworkd || { stop; exit 0; }
    test -r /etc/frameworkd.conf || { stop; exit 0; }
end script

exec /usr/bin/fso-frameworkd
