# fso-deviced.conf
description "fso-deviced - freesmartphone.org device daemon"

start on started dbus
stop on stopping dbus

respawn

pre-start script
    test -x /usr/sbin/fsodeviced || { stop; exit 0; }
end script

exec /usr/sbin/fsodeviced
