# fwlogwatch.conf
description "fwlogwatch - firewall log analyzer"

start on filesystem
stop on runlevel [016]

expect daemon
respawn

pre-start script
    test -x /usr/sbin/fwlogwatch || { stop; exit 0; }
end script

exec /usr/sbin/fwlogwatch -R
