# gdm3.conf
# GDM - Display manager for the GNOME desktop environment
#
# The display manager service manages the X servers running on the
# system, providing login and auto-login services

description "GNOME Display Manager"
author "Cameron Norman <camerontnorman@gmail.com"

start on (filesystem
          and started dbus
          and (plymouth-ready
               or started null-splash))

emits login-session-start
emits desktop-session-start

task

env FORCE_START="no"

pre-start script
    [ -n "$UPSTART_EVENTS" ] || FORCE_START="yes"

    # update the dconf profile if needed
    needed=no
    mkdir -p /etc/dconf/db/gdm.d/
    if [ ! -f /etc/dconf/db/gdm ]; then
	needed=yes
    else
	for f in \
	    /etc/dconf/db/gdm.d /etc/dconf/db/gdm.d/*-settings \
	    /etc/dconf/db/gdm.d/locks /etc/dconf/db/gdm.d/locks/*-locks \
	    /etc/gdm3/greeter.gsettings; do
	    if [ "$f" -nt /etc/dconf/db/gdm ]; then  # newer?
		needed=yes
		break
	    fi
	done
    fi
    if [ "$needed" = yes ]; then
	# The configuration file in /etc uses org.gnome.blah gsettings syntax.
	# Convert it to org/gnome/blah dconf syntax, if needed (newer)
	dist_settings=/etc/dconf/db/gdm.d/90-debian-settings 
	if [ ! -f $dist_settings ] || \
	    [ /etc/gdm3/greeter.gsettings -nt $dist_settings ] ; then
	    awk '/\[.*\]/ { gsub("\\.","/"); } ! /^#/ { print;}' \
		/etc/gdm3/greeter.gsettings > $dist_settings
	fi
	dconf update
    fi
end script

exec start display-manager EXEC="gdm3" FORCE="$FORCE_START"
