# gdomap.conf
description "gdomap - GNUStep distributed object mapper"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

exec gdomap
