# gdomap.conf
description "gdomap - GNUStep distributed object mapper"

start on remote-filesystems and filesystem and started rsyslog
stop on runlevel [016]

expect fork
respawn

exec gdomap
