# getty-ttys.conf
description "Start getty jobs on tty1-6"

start on filesystem and started hostname and started motd and started dbus

env TTYS="1 2 3 4 5 6"

pre-start script
    for n in $TTYS; do
        start -n getty PORT=tty"$n" || continue
    done
end script

post-stop script
    for n in "$TTYS"; do
        stop -n getty PORT=tty"$n" || continue
    done
end script
