# git-daemon.conf
description "git-daemon - simple server for git repositories"

start on filesystem and started networking
stop on deconfiguring-networking or runlevel [016]

expect fork
respawn

env PIDFILE=/var/run/git-daemon.pid

pre-start script
    [ -e /usr/share/git-core/sysvinit/sentinel ] || { stop; exit 0; }
end script

exec /usr/lib/git-core/git-daemon --user=gitdaemon --pid-file=$PIDFILE \
                                  --detach --reuseaddr --verbose

post-stop exec rm -f $PIDFILE
