# gmediaserver.conf
description "GMediaServer - A UPnP compatible media server"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

reload signal 0

setuid nobody
setgid nogroup

pre-start script
    test -x /usr/bin/gmediaserver || { stop; exit 0; }
end script

script
    test ! -s /etc/default/gmediaserver || . /etc/default/gmediaserver
    exec /usr/bin/gmediaserver -b $GMEDIASERVERARGS $GMEDIASERVERDIR
end script
