# gnunet-server.conf
description "GNUnet-server - secure, trust-based peer-to-peer framework"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

setuid gnunet
setgid gnunet

reload signal 0

env ARGS="-c /etc/gnunet.conf -d -l /var/log/gnunetd/gnunetd.log"

# TODO: more than just x86_64 ;)
exec /usr/lib/x86_64-linux-gnu/gnunet/libexec/gnunet-service-arm $ARGS
