# gpsd.conf
description "gpsd - GPS (Global Positioning System) Daemon"

# The GPS daemon should start before NTP, so NTP can get the
# time based on gpsd's services.
start on (filesystem and static-network-up) or starting ntp
stop on deconfiguring-networking or runlevel [016]

expect fork
respawn

env RUNDIR=/run/gpsd

pre-start script
    [ -x /usr/sbin/gpsd ] || { stop; exit 0; }
    [ -d $RUNDIR ] || mkdir -p $RUNDIR
end script

exec /usr/sbin/gpsd -F $RUNDIR/gpsd.sock
