# guacd.conf
description "guacd - Guacamole proxy daemon"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

expect daemon
respawn

setuid guacd
setgid guacd

# No need to create pid dir here, we are not using a pid file
pre-start exec test -x /usr/sbin/guacd || { stop; exit 0; }

exec /usr/sbin/guacd
