# hdapsd.conf
description "hdapsd - park the drive in case of an emergency"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

env DISK=sda

pre-start script
	test -f /sys/block/$DISK/queue/protect -o -f /sys/block/$DISK/device/unload_heads || {
		echo "Not starting, as neither /sys/block/$DISK/queue/protect nor /sys/block/$DISK/device/unload_heads exist."
		echo "Please read /usr/share/doc/hdapsd/README.Debian"
		stop; exit 0
	}
end script

respawn
exec /usr/sbin/hdapsd -d $DISK
