# hostname.sh.conf
description "Set the hostname"

start on startup

pre-start exec hostname -b -F /etc/hostname
