# hostname.conf
description "Set the HOSTNAME environment variable"

start on startup

task

pre-start exec hostname -b -F /etc/hostname
