# hostname.sh.conf
description "Signal sysvinit that HOSTNAME is set"

start on started hostname
