# hwclock.sh.conf
description "Save system clock to hardware clock"

start on startup
stop on runlevel [06]

post-stop exec hwclock --systohc
