# hwclock.sh.conf
description "Save system clock to hardware clock"

start on virtual-filesystems
stop on runlevel [06]

pre-start exec hwclock --hctosys
post-stop exec hwclock --systohc
