# hwclock.sh.conf
description "Relay hwclock.sh readiness for compatibility"

start on started hwclock
stop on runlevel [016] or unmounting-filesystem
