# incron.conf
description "Filesystem event notification daemon"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

respawn
exec /usr/sbin/incrond -n
