# inspircd.conf
description "inspircd - modular C++ irc server"

start on filesystem and static-network-up
stop on runlevel [016]

expect fork
respawn

env INSPIRCD_ENABLED=0
env USER=irc
env LOG="/var/log/inspircd.log"
env CONF="/etc/inspircd/inspircd.conf"
env PIDF="/var/run/inspircd.pid"

pre-start script
    [ -x /usr/sbin/inspircd ] || { stop; exit 0; }
    [ ! -s /etc/default/inspircd ] || . /etc/default/inspircd
    [ "$INSPIRCD_ENABLED" = "1" ] || { stop; exit 0; }
    [ -f $LOG ] || install -m 0640 -o $USER -g adm $LOG
    [ -f $PIDF ] || { touch $PIDF; chown $USER:$USER $PIDF; }
end script

script
    export LD_LIBRARY_PATH=/usr/lib/inspircd
    exec start-stop-daemon -S -q -c $USER -x /usr/sbin/inspircd -- --logfile $LOG --config $CONF
end script
