# ircd-irc2.conf
description "ircd-irc2 - IRC daemon"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start script
    [ -x /usr/sbin/ircd ] || { stop; exit 0; }
    [ -d /var/run/ircd ] || mkdir /var/run/ircd
    chmod 775 /var/run/ircd
end script

exec ircd -c
