# ironic-api.conf
description "Ironic api"
author "Chuck Short <zulcss@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]


chdir /var/run

pre-start script
        mkdir -p /var/run/ironic
        chown ironic:root /var/run/ironic/

        mkdir -p /var/lock/ironic
        chown ironic:root /var/lock/ironic/
end script

exec start-stop-daemon --start --chuid ironic --exec /usr/bin/ironic-api -- --config-file=/etc/ironic/ironic.conf
