# iwatch.conf
description "iwatch - daemon for monitoring the filesystem"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/bin/iwatch || { stop; exit 0; }
end script

exec /usr/bin/iwatch -d -f /etc/iwatch/iwatch.xml
