# kanla.conf
description "kanla - small-scale alerting daemon"

start on filesystem
stop on runlevel [016]

# TODO: readiness
respawn

setuid kanla
setgid kanla

pre-start script
    test -x /usr/bin/kanla || { stop; exit 0; }
end script

exec /usr/bin/kanla
