# kmod.conf
description "Use modall.sh to load kernel modules"

start on startup and filesystem

task

exec /usr/local/bin/modall.sh
