# knockd.conf
description "knockd - port knock daemon"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

umask 0037

pre-start script
    test -x /usr/sbin/knockd || { stop; exit 0; }
end script

exec /usr/sbin/knockd -d
