# kxd.conf
description "kxd - Key exchange daemon"

start on filesystem
stop on runlevel [016]

respawn

pre-start script
    test -x /usr/bin/kxd || { stop; exit 0; }
end script

script
    test ! -r /etc/default/kxd || . /etc/default/kxd
    exec /usr/bin/kxd $OPTS
end script
