# latd.conf
description "latd - Local Area Transport Daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

expect fork
respawn
reload signal 0

env LATCP=/usr/sbin/latcp

pre-start script
    test -x "$LATCP" || { stop; exit 0; }
end script

exec "$LATCP" -s -d

pre-stop exec "$LATCP" -h
