# LCDd.conf
description "LCDd - listen to a port and display info on an LCD display"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

respawn

pre-start script
    test -x /usr/sbin/LCDd || { stop; exit 0; }
end script

exec /usr/sbin/LCDd -f -s 1 -c /etc/LCDd.conf
