# libvirt-bin.conf
description "libvirt-bin - Virtualization daemon"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem

expect daemon
respawn

# Override the maximum number of opened files
#limit nofile 2048 2048

script
    [ ! -s /etc/default/libvirt-bin ] || . /etc/default/libvirt-bin
    exec /usr/sbin/libvirtd -d $libvirtd_opts
end script
