# lightdm.conf
# LightDM - light Display Manager
#
# The display manager service manages the X servers running on the
# system, providing login and auto-login services

description "LightDM Display Manager"
author "Cameron Norman <camerontnorman@gmail.com"

start on (filesystem
          and started dbus
          and plymouth-ready)
stop on (stopped display-manager EXEC="lightdm")

emits login-session-start
emits desktop-session-start

emits stop-display-manager

env FORCE_START="no"

pre-start script
	[ -z "$UPSTART_EVENTS" ] && FORCE_START="yes"
	exec start display-manager EXEC="lightdm" FORCE="$FORCE_START"
end script

post-stop exec initctl emit -n stop-display-manager EXEC="lightdm"
