# lighttpd.conf
description "lighttpd - fast and small webserver"

start on filesystem and static-network-up
stop on runlevel [016]

# No readiness :(
respawn

env CONF=/etc/lighttpd/lighttpd.conf
env RUNDIR=/var/run/lighttpd

pre-start script
    [ -x /usr/sbin/lighttpd ] || { stop; exit 0; }
    /usr/sbin/lighttpd -t -f $CONF || { stop; exit 0; }
    [ -d $RUNDIR ] || install -d -o www-data -g www-data -m 0750 $RUNDIR
end script

exec /usr/sbin/lighttpd -D -f $CONF
