# lm-sensors.conf
description "lm-sensors - hardware health monitoring"

start on filesystem
stop on runlevel [016]

pre-start exec sensors -s
