# lnpd.conf
description "lnpd - BrickOS LNP daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

# TODO: readiness
respawn
reload signal 0

pre-start script
    test -x /usr/sbin/lnpd || { stop; exit 0; }
end script

script
    if test -r /etc/lnpd.conf; then
        . /etc/lnpd.conf
    fi

    exec /usr/sbin/lnpd -d $OPT_LOG $OPT_DEBUG $OPT_FAST $OPT_TTY
end script
