# loadkeys.conf
start on virtual-filesystems or starting mountall-shell

exec loadkeys /etc/console-setup/cached.kmap.gz
