# lprng.conf
description "lprng - Next generation line printer daemon"

start on runlevel [2345]
stop on runlevel [016] or deconfiguring-networking or unmounting-filesystem

pre-start script
	test -f /etc/printcap || test -f /etc/lprng/printcap || {
		echo "No file /etc/printcap or /etc/lprng/printcap exists, not starting."
		stop; exit 0
	}
	mkdir -p /var/run/lprng
	checkpc -f || true
end script

respawn
expect fork
exec /usr/sbin/lpd

post-stop exec rm -rf /var/run/lprng /var/spool/lpd/*/lock.pr /var/spool/lpd/*/control.pr /var/spool/lpd/*/unspooler.pr
