# lsh-server.conf
description "lsh-server - secure shell (SSH2) server"

start on filesystem and static-network-up
stop on deconfiguring-networking or runlevel [016]

# TODO: readiness
respawn

# Default, change in /etc/default/lsh-server
env LSHD_PORT=22

pre-start script
    test -x /usr/sbin/lshd || { stop; exit 0; }
end script

script
    [ ! -r /etc/default/lsh-server ] || . /etc/default/lsh-server
    case "$ENABLE_SFTP" in
      true|y*|Y*)
        SFTP="--subsystems sftp=/usr/lib/lsh-server/sftp-server"
        ;;
    esac

    exec /usr/sbin/lshd --port "$LSHD_PORT" $SFTP $EXTRA_FLAGS
end script
