# lsyncd.conf
description "Local file syncing daemon"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

nice 10

respawn
expect fork
exec /usr/bin/lsyncd /etc/lsyncd/lsyncd.conf.lua
