# lttng-sessiond.conf
description "LTTng 2.0 central tracing registry session daemon"
author "Stéphane Graber <stgraber@ubuntu.com>"

start on local-filesystems
stop on runlevel [06]

respawn
exec lttng-sessiond
