# lvm2.conf
description "lvm2 - Linux Logical Volume Manager, 2"

start on started udev
stop on runlevel [06]

pre-start script
    [ -x /sbin/vgchange ] || { stop; exit 0; }
end script

exec /sbin/lvm vgchange -aay --sysinit
