# lwresd.conf
description "Lightweight Resolver Daemon"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

respawn
exec /usr/sbin/lwresd -f
