# mandos.conf
description "mandos - Gives encrypted passwords to authenticated Mandos clients"

start on started avahi-daemon
stop on stopping avahi-daemon

respawn
respawn limit 1 30

exec /usr/sbin/mandos --foreground

reload signal 0
