# masqmail.conf
description "masqmail - offline Mail Transfer Agent"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

pre-start script
	test -x /usr/sbin/masqmail || { stop; exit 0; }
	test -r /etc/default/masqmail && . /etc/default/masqmail
	test "$INIT_SMTP_DAEMON" = true -o "$INIT_QUEUE_DAEMON" = true || { stop; exit 0; }
	exec mkdir -p /var/run/masqmail
end script

respawn
script
	test -r /etc/default/masqmail && . /etc/default/masqmail
	test "$INIT_SMTP_DAEMON" = true && OPTS="$OPTS -bd"
	test "$INIT_QUEUE_DAEMON" = true && OPTS="$OPTS ${QUEUE_DAEMON_IVAL:--q10m}"
	exec /usr/sbin/masqmail $OPTS
end script
