# memcached.conf
description "memcached - high-performance memory caching daemon"

start on filesystem
stop on runlevel [016]

expect fork
respawn

# Instance w/ default value so main server starts on boot
env SERVER=""
instance $SERVER
usage "SERVER - instance of memcached with corresponding /etc/memcached/\$SERVER config"

# Override these in the instance or global config file to adjust them
env PORT="11211"
env ADDR="127.0.0.1"
env USER="nobody"

pre-start script
    [ -x /usr/bin/memcached ] || { stop; exit 0; }
    if [ "x$SERVER" != "x"]; then
        # Instances need to have their own separate config, or they will
        # try to bind to the same port as the main server and fail.
        [ -s /etc/memcached/$SERVER ] || { stop; exit 0; }
    fi
end script

script
    if [ -s /etc/memcached/global ]; then . /etc/memcached/global; fi
    # If we are an instance, need to source the config
    if [ "x$SERVER" != "x"]; then . /etc/memcached/$SERVER; fi

    exec /usr/bin/memcached -d -p $PORT -l $ADDR -u $USER $OPTIONS
end script
