# memlockd.conf
description "memlockd - locks files in memory to mitigate the effects of a DoS attack"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/sbin/memlockd || { stop; exit 0; }
end script

exec /usr/sbin/memlockd -u memlockd
