# minissdpd.conf
description "MiniSSDPd - UPnP devices daemon"

start on filesystem
stop on runlevel [016]

expect fork
respawn

env START_DAEMON=1
env MiniSSDPd_INTERFACE_ADDRESS=0.0.0.0

pre-start script
    [ ! -s /etc/default/minissdpd ] || . /etc/default/minissdpd
    [ "${START_DAEMON}" = "1" ] || { stop; exit 0; }
    [ -n ${MiniSSDPd_INTERFACE_ADDRESS} ] || { stop; exit 0; }
end script

script
    [ ! -s /etc/default/minissdpd ] || . /etc/default/minissdpd
    exec /usr/sbin/minissdpd -i ${MiniSSDPd_INTERFACE_ADDRESS}
end script

# Upstart does not use PID files, but we do not want
# this to be stale anyway.
post-stop exec rm -f /var/run/minissdpd.pid