# miredo.conf
description "miredo - Teredo IPv6 tunneling daemon"

start on filesystem and static-network-up
stop on runlevel [016]

# TODO: readiness
respawn

env RUNDIR=/var/run/miredo
env CONF=/etc/miredo/miredo.conf

pre-start script
    [ -x /usr/sbin/miredo ] || { stop; exit 0; }
    [ -d $RUNDIR ] || mkdir -p -m 0755 $RUNDIR
    /usr/sbin/miredo-checkconf -f $CONF || { stop; exit 0; }
end script

exec /usr/sbin/miredo -f -c $CONF
