# miredo-server.conf
description "miredo-server - Teredo IPv6 server daemon"

start on filesystem and static-network-up
stop on runlevel [016]

respawn

env RUNDIR=/var/run/miredo-server
env CONF=/etc/miredo/miredo-server.conf
env USER=miredo-server

pre-start script
    [ -x /usr/sbin/miredo-server ] || { stop; exit 0; }
    [ -d $RUNDIR ] || install -d -o $USER -g $USER -m 0755 $RUNDIR
end script

exec /usr/sbin/miredo-server -f -c $CONF -u $USER
