# modemmanager.conf
description "Manages modem connections for NetworkManager"

start on starting network-manager
stop on stopped network-manager or runlevel [016]

respawn

exec modem-manager --log-level=ERR
