# monit.conf
description "monit - service and resource monitoring daemon"

start on filesystem and started networking
stop on runlevel [016] or unmounting-filesystem

expect daemon
respawn

env CONFIG=/etc/monit/monitrc

pre-start script
    [ -x /usr/bin/monit ] || { stop; exit 0; }
    [ -s /etc/default/monit ] && . /etc/default/monit
    [ "$START" = "yes" ] || { stop; exit 0; }
end script

script
    [ -s /etc/default/monit ] && . /etc/default/monit
    exec /usr/bin/monit -c $CONFIG $MONIT_OPTS
end script
