# mouseemu.conf
description "mouseemu - mouse emulation daemon"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

respawn

pre-start script
	modprobe -q uinput || :
	# Give udev some time to create the device node
	sleep 1
end script

script
	test -f /etc/default/mouseemu && {
		. /etc/default/mouseemu
		OPTS="$MID_CLICK $RIGHT_CLICK $SCROLL $TYPING_BLOCK"
	}

	exec /usr/sbin/mouseemu -nofork $OPTS
end script
