# mumble-server.conf
description "mumble-server - Murmur VoIP server"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect daemon
respawn

# This controls how many file descriptors the murmur process can open.
# As a rule of thumb, you should have about 20 descriptors per virtaul 
# server and one for each client. So 30 servers with 20 clients each would
# need at least 720 descriptors (30 * 4 + 30 * 20).
#limit nofile $softlimit $hardlimit

env CONF=/etc/mumble-server.ini

env RUNDIR=/var/run/mumble-server
env RUNUSER=mumble-server
env RUNGROUP=mumble-server

pre-start script
    [ -x /usr/sbin/murmurd ] || { stop; exit 0; }
    [ -d $RUNDIR ] || mkdir -p $RUNDIR
    chown $RUNUSER:$RUNGROUP $RUNDIR
end script

exec /usr/sbin/murmurd -ini $CONF
