# mxallowd.conf
description "mxallowd - dynamically whitelist your Mail eXchanger"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016]

expect fork
respawn
reload signal 0

pre-start script
    test -x /usr/sbin/mxallowd || { stop; exit 0; }
end script

exec /usr/sbin/mxallowd
