# mythbuntu-bare-console.conf
# Mythbuntu Bare Console service

description     "Mythbuntu Bare Server - Server"
author          "Thomas Mashos <tgm4883@ubuntu.com>"

start on started networking
stop on stopped networking

exec /usr/share/mythbuntu-bare/bareserver/bareserver.py
respawn
