# nbd-server.conf
description "nbd-server - Network Block Device server"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start exec test -x /bin/nbd-server || { stop; exit 0; }
exec /bin/nbd-server
