# ndpmon.conf
description "ndpmon - IPv6 Neighbor Discovery Protocol Monitor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /usr/sbin/ndpmon
