# network-manager.conf
description "NetworkManager - daemon for automatically connecting to "

start on local-filesystems and started dbus and static-network-up
stop on stopping dbus or runlevel [016]

expect fork
respawn

script
    if [ -r /etc/default/locale ] ; then
        . /etc/default/locale
        export LANG LANGUAGE LC_MESSAGEs LC_ALL
    fi
    
    exec NetworkManager
end script
