# ngircd.conf
description "ngircd - Next generation IRC daemon"

start on filesystem and static-network-up
stop on deconfiguring-networking

expect fork
respawn

env RUNDIR=/var/run/ngircd
env RUNUSER=irc

pre-start script
    [ -x /usr/sbin/ngircd ] || { stop; exit 0; }

    [ -d $RUNDIR ] || mkdir -p $RUNDIR
    chown $RUNUSER:$RUNUSER $RUNDIR

    home=$(getent passwd "$RUNUSER" | cut -d: -f6)
    [ -d $home ] || mkdir -p $home
    chown $RUNUSER:$RUNUSER $home
end script

exec /usr/sbin/ngircd

post-stop exec rm -f $RUNDIR/ngircd.pid
