# ntp.conf
description "NTP - Network Time Protocol daemon"

start on filesystem and static-network-up
stop on deconfiguring-networking

expect fork
respawn

pre-start exec test -x /usr/sbin/ntpd || { stop; exit 0; }

script
    [ ! -e /var/lib/ntp/ntp.conf.dhcp ] || CONF="-c /var/lib/ntp/ntp.conf.dhcp"

    exec /usr/sbin/ntpd -g -u ntp:ntp -p /run/ntpd.pid
end script

post-stop exec rm -f /run/ntpd.pid
