# ntp.conf
description "NTP - Network Time Protocol daemon"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

respawn

script
    [ ! -e /var/lib/ntp/ntp.conf.dhcp ] || CONF="-c /var/lib/ntp/ntp.conf.dhcp"
    exec /usr/sbin/ntpd -n -g -u ntp:ntp $CONF
end script
