# null-splash-shutdown.conf
# plymouth - Userspace bootsplash utility
#
# plymouth provides a boot splash screen on the system console using
# the kernel framebuffer device.  On boot, this is nominally started by
# the initramfs so the pre-start script, script and post-start script
# parts are actually not run.  These are normally run on shutdown instead.

description	"Fake bootsplash utility"

start on runlevel [016]
stop on stopped rc RUNLEVEL=[016]

# Keep the following in sync with null-splash
pre-start script
    # Check if plymouth is installed and enabled
    if [ -x /sbin/plymouthd ] ||
       [ -x /bin/plymouth ] ||
       [ -e /etc/init/plymouth.conf ] ||
       [ -x /etc/init.d/plymouth ]; then
        stop
        exit 0
    fi

    have_manual = "no"
    for f in /etc/init/plymouth.conf /etc/init/plymouth.conf.override; do
        [ -r "$f" ] || continue
        for line in $(cut -d " " -f 1 "$f"); do
            [ "$line" = "manual" ] && have_manual = "yes"
        done
    done
    if [ "$have_manual" != "yes" ]; then
        stop
        exit 0
    fi
end script
