# null-splash.conf
# null-splash - A fake splash with the same job semantics as plymouth
#
# This job is used by display managers to know when to start when
# there is no splash in use. It will automatically fail if plymouth
# is installed and enabled.

description "null-splash - fake splash with the job semantics of plymouth"

start on (startup and
              (graphics-device-added PRIMARY_DEVICE_FOR_DISPLAY=1
               or drm-device-added PRIMARY_DEVICE_FOR_DISPLAY=1
               or stopped udev-fallback-graphics))

stop on (starting display-manager
         or starting ubiquity
         or starting oem-config
         or starting uxlaunch
         or stopped rc RUNLEVEL=[2345]
         or starting mountall-shell)

# Keep the following in sync with null-splash-shutdown
pre-start script
    # Check if plymouth is installed and enabled
    if [ -x /sbin/plymouthd ] ||
       [ -x /bin/plymouth ] ||
       [ -e /etc/init/plymouth.conf ] ||
       [ -x /etc/init.d/plymouth ]; then
        stop
        exit 0
    fi

    have_manual = "no"
    for f in /etc/init/plymouth.conf /etc/init/plymouth.conf.override; do
        [ -r "$f" ] || continue
        for line in $(cut -d " " -f 1 "$f"); do
            [ "$line" = "manual" ] && have_manual = "yes"
        done
    done
    if [ "$have_manual" != "yes" ]; then
        stop
        exit 0
    fi
end script
