# nvtv.conf
description "nvtv — backend server to control TV output hardware"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start script
    test -x /usr/sbin/nvtvd || { stop; exit 0; }
end script

exec /usr/sbin/nvtvd
