# ofono.conf
description "oFono - Mobile Telephony Stack"

start on started dbus
stop on stopping dbus

expect fork
respawn

exec ofonod
