# olsrd.conf
description "olsrd - olsrd - Optimized Link State Routing protocol daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

# Default value as fallback, change in /etc/default/olsrd
env DAEMON_OPTS="-i wlan0"

pre-start script
    test -x /usr/sbin/olsrd || { stop; exit 0; }
    test ! -s /etc/default/olsrd || . /etc/default/olsrd
    test "$START_OLSRD" = "YES" || { stop; exit 0; }
    if test "$SETUP_ADHOC" = "YES" ; then
        /usr/sbin/olsrd-adhoc-setup $MESH_IF $channel $ssid $bssid
    fi
end script

script
    test ! -s /etc/default/olsrd || . /etc/default/olsrd
    exec /usr/sbin/olsrd $DAEMON_OPTS
end script
