# onak.conf
description "onak - PGP keyserver backend"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/sbin/keyd || { stop; exit 0; }
    egrep -q '^use_keyd *(true|yes|1)$' /etc/onak.conf || {
        echo "Disabled, so not starting (enable in /etc/onak.conf)."
        stop; exit 0
    }
end script

exec start-stop-daemon -S -q -c onak -x /usr/sbin/keyd
