# openbsd-inetd.conf
description "openbsd-inetd - internet 'super-server'"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start script
    [ -x /usr/sbin/inetd ] && [ -e /etc/inetd.conf ] || { stop; exit 0; }
    # Services must be enabled if this is to be started
    grep -q "^[[:alnum:]/]" /etc/inetd.conf || { stop; exit 0; }
end script

exec /usr/sbin/inetd
