# openhpid.conf
description "openhpid - HPI instance to which multiple clients can connect"

start on runlevel [2345]
stop on runlevel [016] or deconfiguring-networking or unmounting-filesystem

exec /usr/sbin/openhpid -n -c /etc/openhpi/openhpi.conf
