# openvpn.conf
description "Automatically spawns OpenVPN instances, if configured"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

# Default in case /etc/default/openvpn is not available.
# Do not change here. Change the value there to configure this job.
env AUTOSTART="all"

pre-start script
    [ -x /usr/sbin/openvpn ] || { stop; exit 0; }
    [ -d /run/openvpn ] || mkdir -p /run/openvpn
end script

script
    [ -r /etc/default/openvpn ] && . /etc/default/openvpn

    # If $AUTOSTART is blanked, interpret that as "all"
    [ "x$AUTOSTART" != "x" ] || AUTOSTART="all"

    if [ "$AUTOSTART" = "none" ]; then
        # We start nothing
        exit 0
    elif [ "$AUTOSTART" = "all" ]; then
        # Start everything in the config directory
        for conf in /etc/openvpn/*.conf; do
            vpn=$(basename ${conf%%.conf})
            # Use the no wait option so that new spawns happen before the
            # first one is finished (faster startup). The downside is that
            # this job can not signal that all VPNs are up.
            start --no-wait openvpn VPN=$vpn
        done
    else
        # $AUTOSTART is a list of VPNs, start only those
        for vpn in $AUTOSTART; do
            # The instance job will check if it has a config file;
            # no need to check that here.
            start --no-wait openvpn VPN=$vpn
        done
    fi
end script

post-stop exec initctl emit stop-openvpns
