# opensm.conf
description "opensm - subnet manager"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

pre-start script
	test -f /sys/class/infiniband_mad/abi_version || {
		echo "No infiniband adapters found."
		exit 0
	}

	if test -f /etc/default/opensm; then
		. /etc/default/opensm
	fi

	test "$PORTS" = "ALL" && PORTS=$(/usr/sbin/ibstat -p)
	test "$PORTS" = "NONE" && exit 0

	for p in $PORTS; do
		start opensm-instance PORT=$p || err=y
	done

	test -n "$err" && exit 1
end script
