# opendkim.conf
description "OpenDKIM - DKIM signing and verifying filter for MTAs"

# TODO: named
start on filesystem and net-device-up IFACE=lo
stop on runlevel [016]

respawn
reload signal USR1

env USER=opendkim
env RUNDIR=/run/opendkim
env CONF=/etc/opendkim.conf
# Modify this one in the default file (if available)
env SOCKET="local:/run/opendkim/opendkim.sock"

pre-start script
    test -x /usr/sbin/opendkim || { stop; exit 0; }

    test -d $RUNDIR || {
        install -o $USER -g $USER -m 755 -d $RUNDIR
        test ! -x /sbin/restorecon || /sbin/restorecon "$RUNDIR"
    }
end script

script
    # Get socket from config
    SOCKET=$(awk '$1 == "Socket" { print $2 }' $CONF)

    # Get socket from default file (override config location)
    if test -r /etc/default/opendkim ; then
        . /etc/default/opendkim
    fi

    # If either of those were invalid, reset SOCKET to the default
    test -n "$SOCKET" || SOCKET="local:/run/opendkim/opendkim.sock"

    # Strip local: prefix to get path name
    SOCK_FILE="${SOCKET#local:}"

    # Remove leftover socket file, if there is one
    if test -f "$SOCK_FILE" -a -S "$SOCK_FILE" ; then
        rm -f "$SOCK_FILE"
    fi

    exec /usr/sbin/opendkim -f -x ${CONF} -u ${USER} -p ${SOCKET} $DAEMON_OPTS
end script

post-stop exec rm -rf $RUNDIR
