# openjpip.conf
description "openjpip - JPIP server for OpenJPEG"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/lib/openjpip-server/opj_server || { stop; exit 0; }
end script

exec spawn-fcgi -d /var/lib/openjpip-server -u www-data -g www-data -p 3000 -- /usr/lib/openjpip-server/opj_server

pre-stop exec GET "http://localhost/openjpip?quitJPIP"
