# papercut.conf
description "papercut - NNTP server"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and net-device-up IFACE=lo
stop on runlevel [016]

expect daemon
respawn

setuid papercut
setgid papercut

pre-start script
    test -x /usr/bin/papercut || { stop; exit 0; }
    # Daemon gets angry if there is already a pid file.
    # This causes respawn and restart to fail.
    rm -f /var/run/papercut/papercut.pid
end script

exec /usr/bin/papercut

post-stop exec rm -f /var/run/papercut/papercut.pid
