# partimaged.conf
description "partimaged - backup/restore partitions server"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

respawn
reload signal 0

env TARGET=/var/lib/partimaged/

pre-start script
    test -x /usr/sbin/partimaged || { stop; exit 0; }
end script

exec /usr/sbin/partimaged -d $TARGET
