# patchram.conf
# patchram -- Broadcom brcm_patchram_plus utility.

description     "bluetooth initialization"

start on starting bluetooth
stop on stopping bluetooth

expect fork

env PATCHRAM_ARGS="--patchram /lib/firmware/bcm4330.hcd
	--no2bytes \
	--scopcm=0,2,0,0,0,0,0,0,0,0 \
	--enable_hci \
	--enable_lpm \
	--baudrate 3000000 --use_baudrate_for_download \
	--tosleep=50000 \
	/dev/ttyHS2"

exec /usr/bin/brcm_patchram_plus $PATCHRAM_ARGS &

pre-start script
	rfkill unblock bluetooth
	chmod 0660 /sys/class/rfkill/rfkill0/state
	chmod 0660 /sys/class/rfkill/rfkill0/type
	chgrp dialout /sys/class/rfkill/rfkill0/state
	chgrp dialout /sys/class/rfkill/rfkill0/type
end script
