# phonefsod.conf
description "phonefsod - SHR Phone freesmartphone.org daemon"

start on started dbus
stop on stopping dbus

expect fork
respawn

# Daemon can not reload
reload signal 0
# For some reason, daemon can not be stopped properly
kill signal KILL

pre-start exec test -x /usr/sbin/phonefsod || { stop; exit 0; }
exec /usr/sbin/phonefsod
