# pidentd.conf
description "pidentd: implements IDENT user identification protocol"

start on virtual-filesystems

apparmor load /etc/apparmor.d/usr.sbin.identd

# do prep work prior to inet activation
pre-start exec install -m 0755 -o identd -g nogroup -d /run/identd
