# pidentd.conf
# The daemon is inet activated. This job is only used for the prep work (run directory).

description "pidentd - implements IDENT user identification protocol"

start on virtual-filesystems
stop on runlevel [016] or unmounting-filesystem

apparmor load /etc/apparmor.d/usr.sbin.identd

env RUNDIR=/run/identd

pre-start script
    [ -x /usr/sbin/identd ] || { stop; exit 0; }
    [ -d $RUNDIR ] || mkdir -p $RUNDIR
    chown identd:nogroup $RUNDIR
    chmod 755 $RUNDIR
end script
