# plymouth-initrd.conf
description "used to tell whether plymouth was started by initrd"

start on startup
stop on started plymouth-stop

# Stop this fucker in pre-start if plymouth is not running so
# that the started event is never emitted for this job.
pre-start exec plymouth --ping || { stop; exit 0; }
