# plymouth-ready.conf
# plymouth-ready - send an event to indicate plymouth is up
#
# This job notifies that the plymouth splash screen is up, either started
# from the initrd or when the 'plymouth-splash' job starts.  The login
# managers use the event to avoid a race with plymouth-splash:
#
# https://bugs.launchpad.net/ubuntu/+source/plymouth/+bug/982889

description "Send an event to indicate plymouth is up"

start on ((started plymouth-initrd or started plymouth-splash)
          and (graphics-device-added PRIMARY_DEVICE_FOR_DISPLAY=1
               or drm-device-added PRIMARY_DEVICE_FOR_DISPLAY=1
               or stopped udev-fallback-graphics))

emits plymouth-ready

task

exec initctl emit plymouth-ready
