# polipo.conf
description "polipo - caching web proxy"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

setuid proxy
setgid proxy

env CONF=/etc/polipo/config

respawn
exec /usr/bin/polipo -c $CONF daemonise=false

reload signal USR2
