# pommed.conf
description "pommed - Apple hotkeys event handler"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

respawn

pre-start exec test -x /usr/sbin/pommed || { stop; exit 0; }
exec /usr/sbin/pommed -f
