# powerwaked.conf
# powerwaked

description "PowerWake Server"
author "Andres Rodriguez <andreserl@canonical.com>"

start on runlevel [2345]

expect daemon

exec /usr/sbin/powerwaked
