# pptpd.conf
description "pptpd - Point to Point Tunneling Server"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start script
    test -x /usr/sbin/pptpd || { stop; exit 0; }
end script

exec /usr/sbin/pptpd

post-stop exec rm -f /var/run/pptpd.pid
