# preload.conf
description "preload - adaptive readahead daemon"

start on filesystem
stop on runlevel [016]

expect fork
respawn

env DAEMON_OPTS="-s /var/lib/preload/preload.state"

pre-start script
    [ -x /usr/sbin/preload ] || { stop; exit 0; }
end script

script
    [ -s /etc/default/preload ] && . /etc/default/preload
    DAEMON_OPTS="$DAEMON_OPTS $OPTIONS"
    # IO scheduling value not supported because it is too difficult to do

    exec preload $DAEMON_OPTS
end script
