# privoxy.conf
description "privoxy - Privacy enhancing proxy"

start on filesystem and static-network-up
stop on deconfiguring-networking or runlevel [016]

expect fork
respawn

env CONF=/etc/privoxy/config
env USER=privoxy
env LOGDIR=/var/log/privoxy

pre-start script
    [ -x /usr/sbin/privoxy ] || { stop; exit 0; }
    [ -d $LOGDIR ] || install -d -o $USER -g adm -m 750 $LOGDIR
end script

exec /usr/sbin/privoxy --user $USER $CONF
