# proftpd.conf
description "proftpd - Pro FTP daemon"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

expect fork
respawn

# Change in the default file, not here
env RUN=yes

pre-start script
    [ -x /usr/sbin/proftpd ] || { stop; exit 0; }
    [ ! -s /etc/default/proftpd ] || . /etc/default/proftpd
    [ "x$RUN" = "xyes" ] || { stop; exit 0; }
    [ -d /run/proftpd ] || mkdir -p /run/proftpd

    # Check for inetd
    if egrep -qi "server[[:space:]]*=[[:space:]]*/usr/sbin/(in\.)?proftpd" /etc/xinetd.conf 2>/dev/null || \
       egrep -qi "server[[:space:]]*=[[:space:]]*/usr/sbin/(in\.)?proftpd" /etc/xinetd.d/* 2>/dev/null || \
       egrep -qi "^ftp.*/usr/sbin/(in\.)?proftpd" /etc/inetd.d/* 2>/dev/null || \
       egrep -qi "^ftp.*/usr/sbin/(in\.)?proftpd" /etc/inetd.conf 2>/dev/null
	then
        echo "Configured via inetd; not starting."
        stop; exit 0
    fi
end script

exec /usr/sbin/proftpd -c /etc/proftpd/proftpd.conf
