# ptpd.conf
description "ptpd - Precision Time Protocol daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

pre-start script
    test -x /usr/sbin/ptpd || { stop; exit 0; }
end script

script
    if test -r /etc/default/ptpd; then
        . /etc/default/ptpd
    fi
    
    exec /usr/sbin/ptpd $PTPD_OPTS
end script
