# pyro-nsd.conf
description "pyro-nsd - Pyro name server daemon"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

expect daemon
respawn

pre-start script
    test -x /usr/bin/pyro-nsd || { stop; exit 0; }
end script

exec /usr/bin/pyro-nsd start

pre-stop exec /usr/bin/pyro-nsd stop
