# pyroman.conf
description "pyroman - firewall"

start on (starting network-interface or
          starting networking)

pre-start script
    test -x /usr/sbin/pyroman || { stop; exit 0; }
    /usr/sbin/pyroman --init
end script
