# qpidd.conf
description "qpidd - AMQP message broker"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

respawn
setuid qpidd
setgid qpidd

env SASL_CONF_PATH=/etc/sasl2

pre-start script
    test -x /usr/sbin/qpidd || { stop; exit 0; }
end script

exec /usr/sbin/qpidd --config /etc/qpid/qpidd.conf
