# qtsmbstatusd.conf
description "qtsmbstatusd - Graphical monitor of samba"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and started smbd
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/bin/qtsmbstatusd || { stop; exit 0; }
end script

exec /usr/bin/qtsmbstatusd -d
