# quake3-server.conf
description "quake3-server - server for Quake, a first person shooter game"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

respawn

setuid Debian-quake3
setgid games

pre-start exec test -x /usr/lib/ioquake3/ioq3ded || { stop; exit 0; }
exec /usr/lib/ioquake3/ioq3ded +exec debian_server.cfg
