# quake2-server.conf
description "quake2-server - server for Quake, a first person shooter game"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

respawn

setuid quake2-server
setgid games

pre-start exec test -x /usr/lib/quake2/quake2-engine-server || { stop; exit 0; }
exec /usr/lib/quake2/quake2-engine-server +exec debian_server.cfg
