# quake-server.conf
description "quake-server - server for Quake, a first person shooter game"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

respawn

setuid quake-server
setgid games

pre-start exec test -x /usr/lib/quake/quake-engine-server || { stop; exit 0; }
exec /usr/lib/quake/quake-engine-server +exec debian_server.cfg
