# radicale.conf
description "radicale - a simple calendar and contact server"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

# TODO: readiness
respawn

umask 0027
setuid radicale
setgid radicale

pre-start script
    test -x /usr/bin/radicale || { stop; exit 0; }
end script

exec /usr/bin/radicale -f
reload signal 0
