# radsecproxy.conf
description "radsecproxy - RADIUS protocol proxy supporting RadSec"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and net-device-up IFACE=lo
stop on runlevel [016]

expect fork
respawn

pre-start exec test -x /usr/sbin/radsecproxy || { stop; exit 0; }
exec /usr/sbin/radsecproxy
