# ramond.conf
description "ramond - Router Advertisement MONitoring Daemon"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

expect daemon
respawn

env CONF=/etc/ramond.conf

# Daemon can not reload via signal
reload signal 0

pre-start script
    test -x /usr/sbin/ramond || { stop; exit 0; }
    test -f $CONF
end script

exec /usr/sbin/ramond -c $CONF
