# rarpd.conf
description "rarpd - Reverse Address Resolution Protocol (RARP) daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start script
    test -x /usr/sbin/rarpd || { stop; exit 0; }
end script

exec /usr/sbin/rarpd -e -v
