# reconserver.conf
description "reConServer - SIP comferencing server"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

umask 002

pre-start exec install -o recon -g recon -d /run/reConServer /var/log/reConServer

exec /sbin/start-stop-daemon -S -q -c recon:recon -x /usr/sbin/reConServer -- \
	/etc/reConServer/reConServer.config --Daemonize=false --PidFile=/run/reConServer/reConServer.pid
