# redis-server.conf
description "redis-server - Persistent key-value database"

start on filesystem
stop on runlevel [016]

expect fork

env CONF=/etc/redis/redis.conf
env RUNDIR=/var/run/redis

umask 007

pre-start script
    [ -x /usr/bin/redis-server ] || { stop; exit 0; }
    
    [ -d $RUNDIR ] || mkdir -p $RUNDIR
    touch $RUNDIR/redis-server.pid
    chown redis:redis $RUNDIR $RUNDIR/redis-server.pid
    chmod 755 $RUNDIR
end script

script
    [ -s /etc/default/redis-server ] && . /etc/default/redis-server
    [ -n $ULIMIT ] && ulimit -n $ULIMIT

    exec redis-server $CONF
end script

post-stop exec rm -f $RUNDIR/redis-server.pid
