# redsocks.conf
description "TCP connection redirector to a proxy server"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

respawn
expect fork
exec /usr/sbin/redsocks -c /etc/redsocks.conf
