# restartd.conf
description "restard - process checker"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/sbin/restartd || { stop; exit 0; }
end script

exec /usr/sbin/restartd
