# rfkill.conf
description "rfkill: persist software rfkill state across reboots"

start on local-filesystems or unmounting-filesystem

env UPSTART_EVENTS="local-filesystems"
instance $UPSTART_EVENTS

task

script
    test -d /sys/class/rfkill ] || exit 0
    action="load"
    test "$UPSTART_EVENTS" = "local-filesystems" || action="save"
    for device in /sys/class/rfkill/*; do
        test -d $device || continue
        /lib/systemd/systemd-rfkill $action `basename $device`
    done
end script
