# rhnsd.conf
description "rhnsd - Red Hat Network Daemon"

start on filesystem and net-device-up IFACE!=lo
stop on runlevel [016] or deconfiguring-networking

expect daemon
respawn

pre-start script
    test -x /usr/sbin/rhnsd || { stop; exit 0; }
end script
exec /usr/sbin/rhnsd
post-stop exec rm -f /var/run/rhnsd.pid
