# rrdcached.conf
description "rrdcached - Data caching daemon for rrdtool"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

env DIR="/var/lib/rrdcached/db/"
env JOURNAL="/var/lib/rrdcached/journal/"
env SOCKET="unix:/var/run/rrdcached.sock"

pre-start exec mkdir -p $DIR $JOURNAL

exec /usr/bin/rrdcached -b $DIR -B -j $JOURNAL -F -l $SOCKET
