# rsync.conf
description "rsync - Fast remote file copy daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect daemon
respawn

# The default file should override this, if it is installed
env RSYNC_ENABLE="true"

pre-start script
    [ -s /etc/default/rsync ] && . /etc/default/rsync
    [ "$RSYNC_ENABLE" = "true" ] || { stop; exit 0; }
    [ -s /etc/rsyncd.conf ] || { stop; exit 0; }
end script

exec rsync --daemon --no-detatch
