# rwhod.conf
description "rwhod - system status server"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn
exec /usr/sbin/rwhod -b
