# samhain.conf
description "samhain - check file integrity"

start on filesystem
stop on runlevel [016]

pre-start script
	install -d /var/run/samhain
	test -f /var/state/samhain/samhain_file || {
		echo "Creating integrity database (this can take a few minutes)."
		samhain -t init >/var/log/samhain/samhain-init.log 2>&1
		test -f /var/state/samhain/samhain_file || {
			echo "Database could not be created. Review /var/log/samhain/samhain-init.log"
			exit 1
		}
	}
end script

respawn
exec /usr/sbin/samhain --foreground --forever
