# saned.conf
description "saned - Scanner Access Now Easy"

start on filesystem
stop on runlevel [016]

expect fork
respawn

env RUN=no

setuid saned
setgid saned

script
    [ -s /etc/default/saned ] && . /etc/default/saned
    [ "x$RUN" = "xyes" ] || { stop; exit 0; }
end script
    
exec saned -a
