# scanlogd.conf
description "scanlogd - logs TCP port scans"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

expect fork
respawn

umask 022

pre-start script
    test -x /usr/sbin/scanlogd || { stop; exit 0; }
    install -m 755 -d /var/run/scanlogd
end script

exec /usr/sbin/scanlogd
