# schroot.conf
description "schroot - recover schroot sessions"

start on filesystem and started networking
stop on runlevel [016]

pre-start script
    if [ ! -x /usr/bin/schroot ]; then
        echo "Package not installed"
        stop
        exit 0
    fi

    # Delete /var/lock/schroot if it's not root:root 0700, or if
    # it's not a directory.
    if [ -e /var/lock/schroot ]; then
        if [ "$(stat -c "%u %g %a" /var/lock/schroot)" != "0 0 700" ]; then
            rm -rf /var/lock/schroot
        fi
        if [ ! -d /var/lock/schroot ]; then
            rm -rf /var/lock/schroot
        fi
    fi
    if [ ! -d /var/lock/schroot ]; then
        mkdir -m 0700 /var/lock/schroot
    fi

    # Perform start action
    [ -s /etc/default/schroot ] && . /etc/default/schroot
    
    echo -n "Start: "
    if [ "$START_ACTION" = "end" ] || [ "$SESSIONS_RECOVER" = "end" ] ; then
        echo "ending all sessions."
        for chroot in `schroot --all-sessions --list --quiet`; do
            schroot --chroot=$chroot --end-session || continue
        done
    else
        echo "recovering all sessions."
        for chroot in `schroot --all-sessions --list --quiet`; do
            schroot --chroot=$chroot --recover-session || continue
        done
    fi
end script

post-stop script
    # Perform stop action
    [ -s /etc/default/schroot ] && . /etc/default/schroot

    echo -n "Stop:  "
    if [ "${STOP_ACTION}" = "end" ]; then
        echo "ending all sessions."
        for chroot in `schroot --all-sessions --list --quiet`; do
            schroot --chroot=$chroot --end-session || continue
        done
    else
        echo "performing no action."
    fi
end script
