# scrapyd.conf
description "scrapyd - run Scrapy spiders"

start on filesystem and static-network-up
stop on runlevel [016]

respawn

chdir /var/lib/scrapyd

pre-start exec test -x /usr/bin/scrapyd || { stop; exit 0; }

exec /usr/bin/scrapyd -u scrapy -g nogroup -l /var/log/scrapyd/scrapyd.log
