# sensord.conf
description "sensord - hardware sensor information logging daemon"

start on started lm-sensors
stop on runlevel [016]

expect fork
respawn

script
    [ ! -s /etc/default/sensord ] || . /etc/default/sensord

    # start the daemon with a bunch of optional args
    exec /usr/sbin/sensord ${ALARM_INTERVAL:+  -i $ALARM_INTERVAL} \
                           ${LOG_INTERVAL:+    -l $LOG_INTERVAL} \
                           ${SYSLOG_FACILITY:+ -f $SYSLOG_FACILITY} \
                           ${RRD_INTERVAL:+    -t $RRD_INTERVAL} \
                           ${RRD_FILE:+        -r $RRD_FILE} \
                           ${RRD_LOADAVG:+     -a} \
                           $CONFIG_FILE \
                           $SCAN_CHIPS
end script