# shibd.conf
description "shibd - Shibboleth 2 Daemon"

start on filesytem and static-network-up
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

env OPTS="-f -w 30 -c /etc/shibboleth/shibboleth2.xml"
env RUNDIR="/run/shibboleth"
env USER="_shibd"

pre-start script
    test -x /usr/sbin/shibd || { stop; exit 0; }
    test -d $RUNDIR || mkdir -p $RUNDIR
    chown -Rh $USER:$USER $RUNDIR
    # Make sure we can actually run the daemon as the appropiate user
    if ! start-stop-daemon -S -q -c $USER -x /usr/sbin/shibd -- -t $OPTS;
    then
        echo "Error: shibd test run failed."
        echo "Is shibd's private key owned by the $USER user?"
        exit 1
    fi
end script

exec start-stop-daemon -S -q -c $USER -x /usr/sbin/shibd -- $OPTS
