# shishi-kdc.conf
description "shishid - Kerberos Key Distribution Center (KDC)"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

respawn
reload signal 0

pre-start script
    test -x /usr/sbin/shishid || { stop; exit 0; }
end script

exec /usr/sbin/shishid
