# siproxd.conf
description "siproxd - Proxy for SIP channels of VoIP clients"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect daemon
respawn

# Reload hackily works because Upstart sends a SIGHUP,
# the daemon dies, then Upstart respawns it.

# Change this in the default file, not here
env ENABLED=0

pre-start script
    test -x /usr/sbin/siproxd || { stop; exit 0; }
    if [ -s /etc/default/siproxd ]; then
        . /etc/default/siproxd
    fi
    [ "$ENABLED" != "0" ] || { stop; exit 0; }
end script

exec /usr/sbin/siproxd
