# sipwitch.conf
description "sipwitch - GNU SIP Witch telephony service"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

env DEFAULT=/etc/default/sipwitch
env LOCK=/var/lock/subsys/sipwitch

pre-start script
	test -r $DEFAULT && . $DEFAULT

	test -n "$PLUGINS" || {
		echo "You must set \$PLUGINS in $DEFAULT for sipwitch to start."
		stop; exit 1
	}
end script

respawn
script
	test -r $DEFAULT && . $DEFAULT

	case "$SECURITY" in
		desktop|Desktop|DESKTOP) secopt="--desktop" ;;
		public|Public|PUBLIC) secopt="--public" ;;
	esac

	unset SECURITY

	exec /usr/sbin/sipw -f $secopt
end script

post-start exec touch $LOCK
