# slapd.conf
description "slapd - Standalone LDAP daemon"

start on filesystem and started networking
stop on deconfiguring-networking

expect fork
respawn

# Will not be affected by default file
env RUNDIR=/var/run/slapd
env CONF=/etc/ldap/slapd.d

# Obey the default file
env SLAPD_USER=openldap
env SLAPD_GROUP=openldap
env SLAPD_SERVICES="ldap:/// ldapi:///"

pre-start script
    [ -x /usr/sbin/slapd ] || { stop; exit 0; }
    [ ! -e /etc/ldap/noslapd ] || { stop; exit 0; }
    [ -d $CONF ] || { stop; exit 0; }
    /usr/sbin/slapd -Tt -F $CONF || { stop; exit 0; }

    [ ! -s /etc/default/slapd ] || . /etc/default/slapd
    [ -d $RUNDIR ] || mkdir -p $RUNDIR
    chown -R $SLAPD_USER:$SLAPD_GROUP $RUNDIR
end script

script
    [ ! -s /etc/default/slapd ] || . /etc/default/slapd

    exec /usr/sbin/slapd -F $CONF -u $SLAPD_USER -g $SLAPD_GROUP \
                         -h "$SLAPD_SERVICES" $SLAPD_OPTIONS
end script

post-start script
    # Backward compatibility with OpenLDAP 2.1 client libraries.
    if [ ! -h /var/run/ldapi ] && [ ! -e /var/run/ldapi ] ; then
        ln -s slapd/ldapi /var/run/ldapi
    fi
end script