# slashem-common.conf
description "slashem-common - save file recovery script"

start on filesystem
stop on runlevel [016]

pre-start exec /usr/lib/games/slashem/recover-all
