# slpd.conf
description "slpd - OpenSLP (Service Location Protocol) daemon"

start on filesystem or starting slapd
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/sbin/slpd || { stop; exit 0; }
end script

exec /usr/sbin/slpd
