# slpd.conf
description "slpd - OpenSLP (Service Location Protocol) daemon"

start on filesystem or starting slapd
stop on runlevel [016]

expect fork
respawn

pre-start exec test -x /usr/sbin/slpd || { stop; exit 0; }

exec /usr/sbin/slpd