# sm-archive.conf
description "sm-archive - a sendmail milter to add recipients to messages"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

# Must be in local:pathname or inet:port@address format
env SOCKET="inet:31950@localhost"

pre-start script
    test -x /usr/sbin/sm-archive && test -f /etc/sm-archive.conf || { stop; exit 0; }
    case $SOCKET in
      local:*)
          dir=$(dirname ${SOCKET#local:})
          test -d $dir || { mkdir -p $dir; chown sm-archive: $dir; }
          test ! -x /usr/sbin/postfix || chown postfix: ${SOCKET#local:}
      ;;
    esac
end script

exec /usr/sbin/sm-archive -d 1 -p $SOCKET
