# smartmontools.conf
description "smartmontools - SMART monitoring daemon"

start on filesystem
stop on runlevel [016]

expect daemon
respawn

pre-start script
    test -x /usr/sbin/smartd || { stop; exit 0; }

    # Try to enable smart devices. Cleanly move on to starting the
    # SMART daemon if any type of failure occurs.
    test -x /usr/sbin/smartctl || exit 0
    [ -s /etc/default/smartmontools ] && . /etc/default/smartmontools || exit 0
    for dev in $enable_smart; do
        /usr/sbin/smartctl --quitemode=errorsonly --smart=on $dev || :
    done
end script

exec /usr/sbin/smartd
