# snmptrapd.conf
description "snmptrapd - Simple Network Management Protocol trap daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on started sysinit
stop on stopping sysinit

respawn

env MIBSDIR=/usr/share/snmp/mibs:/usr/share/snmp/mibs/iana:/usr/share/snmp/mibs/ietf:/usr/share/mibs/site:/usr/share/snmp/mibs:/usr/share/mibs/iana:/usr/share/mibs/ietf:/usr/sha
re/mibs/netsnmp

pre-start script
    test -x /usr/sbin/snmptrapd && test -e /etc/snmp/snmptrapd.conf || { stop; exit 0; }
end script

exec /usr/sbin/snmptrapd -Lsd -f
