# snmpd.conf
description "snmpd - Simple Network Management Protocol daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect daemon
respawn

env SNMPDOPTS="-Lsd -Lf /dev/null -u snmp -g snmp -I -smux,mteTrigger,mteTriggerConf"

pre-start script
    test -x /usr/sbin/snmpd || { stop; exit 0; }
    mkdir -p /run/agentx
end script

script    
    OLD_MIBS_DIR="/usr/share/mibs/site:/usr/share/snmp/mibs:/usr/share/mibs/iana:/usr/share/mibs/ietf:/usr/share/mibs/netsnmp"
    MIBS_DIR="/usr/share/snmp/mibs:/usr/share/snmp/mibs/iana:/usr/share/snmp/mibs/ietf"
    MIBDIRS="$MIBS_DIR:$OLD_MIBS_DIR"

    if test -r /etc/default/snmpd; then
        . /etc/default/snmpd
    fi

    exec /usr/sbin/snmpd $SNMPDOPTS
end script
