# snmpd.conf
description "snmpd - Simple Network Management Protocol daemon"
author "Cameron Nemo <camerontnorman@gmail.com>"

start on started sysinit
stop on stopping sysinit

respawn

env MIBS=""
env MIBSDIR=/usr/share/snmp/mibs:/usr/share/snmp/mibs/iana:/usr/share/snmp/mibs/ietf:/usr/share/mibs/site:/usr/share/snmp/mibs:/usr/share/mibs/iana:/usr/share/mibs/ietf:/usr/share/mibs/netsnmp

pre-start script
    test -x /usr/sbin/snmpd && test -e /etc/snmp/snmpd.conf || { stop; exit 0; }
    mkdir -p /run/agentx
end script

exec /usr/sbin/snmpd -Lsd -Lf /dev/null -u Debian-snmp -g Debian-snmp -I -smux,mteTrigger,mteTriggerConf -f
