# snmptrapd.conf
description "snmptrapd - Simple Network Management Protocol trap daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

pre-start script
    test -x /usr/sbin/snmptrapd || { stop; exit 0; }
end script

exec /usr/sbin/snmptrapd -Lsd
