# softflowd.conf
description "softflowd - Flow-based network traffic analyser"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

env CHROOT_DIR=/var/run/softflowd/chroot
# Set the following in the default file, if available
env INTERFACE=any

pre-start script
    test -x /usr/sbin/softflowd || { stop; exit 0; }
    test ! -s /etc/default/softflowd || . /etc/default/softflowd
    if test -z "$INTERFACE" ; then
      echo "Set INTERFACE value in default file!"; stop; exit 0
    fi
    [ -d $CHROOT_DIR ] || mkdir -p $CHROOT_DIR
end script

script
    [ ! -s /etc/default/softflowd ] || . /etc/default/softflowd
    exec /usr/sbin/softflowd -i "$INTERFACE" $OPTIONS
end script

post-stop exec rm -f /var/run/softflowd.pid
