# sogo.conf
description "SOGo - Groupware server"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

env PREFORK=1
env USER=sogo

pre-start script
	test -r /etc/default/sogo && . /etc/default/sogo
	install -o $USER -g $USER -d /run/sogo
end script

script
	test -r /etc/default/sogo && . /etc/default/sogo
	exec start-stop-daemon -S -q -c $USER -x /usr/sbin/sogod -- \
		-WOWorkersCount $PREFORK \
		-WOPidFile /run/sogo/sogo.pid \
		-WOLogFile /var/log/sogo/sogo.log
end script
