# speech-dispatcher.conf
description ""

start on filesystem
stop on runlevel [016]

expect WHAT
respawn

env RUN=no

pre-start script
    [ -x /usr/bin/speech-dispatcher ] || { stop; exit 0; }
    [ -s /etc/default/speech-dispatcher ] && . /etc/default/speech-dispatcher || :
    [ "x$RUN" = "xyes" ] || { stop; exit 0; }
    # Log and run dir stuff
end script

#start setuid speech-dispatcher
exec /usr/bin/speech-dispatcher
