# speechd-up.conf
description "speechd-up - Interface between speakup and speech-dispatcher"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/bin/speechd-up || { stop; exit 0; }
end script

exec /usr/bin/speechd-up -d -l1
