# spice-vdagent.conf
description "spice-vdagent"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
	mkdir -p /run/spice-vdagentd
	modprobe -q -b uinput
	rm -f /run/spice-vdagentd/spice-vdagentd-sock
end script

exec /usr/sbin/spice-vdagentd
