# squid-prefetch.conf
description "prefetch daemon for squid"

start on started squid
stop on stopping squid

exec /usr/sbin/squid-prefetch
