# squidtaild.conf
description "squidtaild — Squid log monitoring program"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

exec /usr/sbin/squidtaild
post-stop exec rm -f /var/run/squidtaild.pid
