# squid.conf
description "squid - "

start on filesystem and static-network-up
stop on runlevel [016]

expect fork
respawn

# Set both to change the user squid is run as.
#setuid root
#setgid root

limit nofile 1024 1024
umask 027

env SPOOLDIR=/var/spool/squid

pre-start script
    [ -x /usr/sbin/squid ] || { stop; exit 0; }
    if [ -d $SPOOLDIR -a ! -d $SPOOLDIR/00 ]; then
        squid -z
    fi
end script

exec /usr/sbin/squid -D -YC
