# sqwebmail.conf
description "Squirrel Webmail Daemon"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

pre-start script
	if test -f /etc/courier/calendarmode && \
	   test -n "$(cat /etc/courier/calendarmode)" && \
	   test ! -d /var/run/courier/calendar
	then
		mkdir --mode=0775 /var/run/courier/calendar
		mkdir --mode=0700 /var/run/courier/calendar/localcache
		mkdir --mode=0750 /var/run/courier/calendar/private
		mkdir /var/run/courier/calendar/public
		chgrp -R daemon /var/run/courier/calendar
	fi
end script

script
	. /etc/courier/sqwebmaild

	DAEMON=/usr/lib/courier/courier/sqwebmaild

	if test "$(cat /etc/courier/calendarmode)" = "net"; then
		DAEMON="/usr/lib/courier/courier/pcpd server"
	fi

	exec $DAEMON
end script
