# sudo.conf
description "sudo - Provide limited superuser priveleges to specific users"

start on filesystem
stop on runlevel [016]

# Ensure that priveleges do not persist across reboots
pre-start exec [ -d /var/lib/sudo ] && find /var/lib/sudo -exec touch -d @0 '{}' \;
