# supervisor.conf
description "supervisor - control and maintain process state"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

expect fork
respawn

pre-start script
    test -x /usr/bin/supervisord || { stop; exit 0; }
end script

exec /usr/bin/supervisord -c /etc/supervisor/supervisord.conf
