# suricata.conf
description "suricata - Intrusion Detection and Prevention System"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem

pre-start script
	if test -r /etc/default/suricata; then
		. /etc/default/suricata
	fi

	case "$LISTENMODE" in
	nfqueue)
		LISTEN_OPT=" -q $NFQUEUE"
		test -e /proc/net/netfilter/nf_queue || {
			echo "ERROR: nfqueue support not found."
			exit 1
		}
		;;
	pcap)
		LISTEN_OPT=" -i $IFACE"
		;;
	af-packet)
		LISTEN_OPT=" --af-packet"
		;;
	*)
		echo "Unsupported listen mode $LISTENMODE, aborting"
		exit 1
		;;
	esac

	initctl set-env LISTEN_OPT="$LISTEN_OPT"
	initctl set-env CONF="$SURCONF"
end script

exec /usr/bin/suricata -c $CONF $LISTEN_OPT
