# sympa-bulk.conf
description "sympa - Mailing List Manager (bulk)"

start on start-sympa
stop on stop-sympa

setuid sympa
setgid sympa

respawn
exec /usr/lib/sympa/bin/bulk.pl -F
