# sysrqd.conf
description "sysrqd"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

pre-start script
	test -f /etc/sysrqd.secret || {
		echo "Not starting, you must set the password in /etc/sysrqd.secret"
		echo "See /usr/share/doc/README.Debian for more information"
		exit 1
	}
end script

expect fork
respawn
exec /usr/sbin/sysrqd
