# sysstat.conf
description "sysstat"

start on filesystem

pre-start exec /usr/lib/sysstat/debian-sa1
