# thin.conf
description "thin - fast and very simple Ruby web server"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

pre-start exec /usr/bin/thin start --all /etc/thin2.1
post-stop exec /usr/bin/thin stop --all /etc/thin2.1
