# thinkfan.conf
description "thinkfan - fan control tool"

start on filesystem
stop on runlevel [016]

expect fork
respawn

exec /usr/sbin/thinkfan -q

post-stop exec rm -f /var/run/thinkfan.pid
