# transmission-daemon.conf
description "transmission-daemon - Bittorrent server"

start on filesystem and static-network-up
# Transmission needs to notify peers before it dips
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

# Play with this for a lower cpu priority
#nice 10

setuid debian-transmission
setgid debian-transmission

env ENABLE_DAEMON=1

pre-start script
    test -x /usr/bin/transmission-daemon || { stop; exit 0; }
    [ ! -s /etc/default/transmission-daemon ] || . /etc/default/transmission-daemon
    [ "$ENABLE_DAEMON" = "1" ] || { stop; exit 0; }
end script

script
    [ ! -s /etc/default/transmission-daemon ] || . /etc/default/transmission-daemon
    exec /usr/bin/transmission-daemon $OPTIONS
end script
