# uhub.conf
description "uhub - a high performance peer-to-peer hub for the ADC network"

# Only enable if you have configured the service
manual

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

limit nofile 65536 65536

env RUNDIR=/var/run/uhub
env LOGFILE=/var/log/uhub/uhub.log

pre-start script
    test -x /usr/bin/uhub || { stop; exit 0; }
    test -d $RUNDIR || mkdir $RUNDIR
end script

exec /usr/bin/uhub -f -l $LOGFILE
