# ulatencyd.conf
description "ulatencyd"

start on started dbus
stop on stopping dbus

expect fork
respawn

exec ulatencyd --daemonize
