# unattended-upgrades.conf
description "unattended-upgrades - Check if unattended upgrades are being applied"

# This is a no-op on startup, but needs some kind of lifetime to play nice with
# the deb helper scripts and startpar.
start on startup
stop on runlevel [016]

env SHUTDOWN_HELPER=/usr/share/unattended-upgrades/unattended-upgrade-shutdown

post-stop script
    PYTHON=python
    [ -x /usr/bin/python3 ] && PYTHON=python3
    $PYTHON $SHUTDOWN_HELPER
end script
