# unbound.conf
description "unbound - Validating, recursive, and caching DNS resolver"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem or deconfiguring-networking

expect fork
respawn

# Defaults, modify in /etc/default/unbound
env UNBOUND_ENABLE=true
env ROOT_TRUST_ANCHOR_UPDATE=false
env ROOT_TRUST_ANCHOR_FILE=/var/lib/unbound/root.key
env RESOLVCONF=false

env CONF=/etc/unbound/unbound.conf

pre-start script
    test -x /usr/sbin/unbound -a -x /usr/sbin/unbound-checkconf || { stop; exit 0; }
    test ! -s /etc/default/unbound || . /etc/default/unbound

    case $UNBOUND_ENABLE in
      true|yes|1) ;;
      *) echo "Disabled in default file."; stop; exit 0 ;;
    esac

    # Setup chroot directory
    CHROOT_DIR=$(awk '{if ($1 ~ "^chroot" && $2 != "\"\"") print $2}' $CONF|sed -e "s#\"##g")
    if [ -d "$CHROOT_DIR" -a "$CHROOT_DIR" != "$(dirname $CONF)" ]; then
        tar --overwrite -cf - $(echo $(dirname $CONF) | sed 's#^/##') | (cd $CHROOT_DIR && tar -xf -)
    fi

    case $ROOT_TRUST_ANCHOR_UPDATE in
      true|yes|1)
        unbound-anchor -a $ROOT_TRUST_ANCHOR_FILE -v 2>&1 | logger -p daemon.info -t unbound-anchor
        chown unbound:unbound $ROOT_TRUST_ANCHOR_FILE
      ;;
    esac
end script

script
    test ! -s /etc/default/unbound || . /etc/default/unbound
    exec /usr/sbin/unbound $DAEMON_OPTS
end script

post-start script
    test ! -s /etc/default/unbound || . /etc/default/unbound

    case $RESOLVCONF in
      yes|1) RESOLVCONF=true ;;
      *) ;;
    esac

    if $RESOLVCONF && [ -x /sbin/resolvconf ]; then
        unbound-checkconf $CHROOT_DIR/$UNBOUND_CONF -o interface | (
            default=yes
            while read interface; do
                default=no
                if [ "x$interface" = x0.0.0.0 -o "x$interface" = x127.0.0.1 ]; then
                    echo "nameserver 127.0.0.1"
                elif [ "x$interface" = x::0 -o "x$interface" = x::1 ]; then
                    echo "nameserver ::1"
                fi
            done
            if [ $default = yes ]; then
                # unbound defaults to listening on localhost
                echo "nameserver 127.0.0.1"
            fi
        ) | /sbin/resolvconf -a lo.unbound
    fi
end script

post-stop script
    if [ -x /sbin/resolvconf ]; then
        /sbin/resolvconf -d lo.unbound
    fi
end script
