# userv.conf
description "userv - User services (security boundary) daemon"

start on filesystem
stop on runlevel [016]

respawn
# Exits with 1 when it gets a SIGTERM/INT
normal exit 0 1
# Cannot reload config in place
reload signal 0

pre-start script
     test -x /usr/sbin/uservd || { stop; exit 0; }
     test -d /var/run/userv || mkdir -p -m 700 /var/run/userv
end script

exec /usr/sbin/uservd
