# virtlockd.conf
description "virtlockd - Virtual machine lock manager"

start on filesystem
stop on runlevel [016]

expect daemon
respawn

oom score -900

reload signal USR1

pre-start exec test -x /usr/sbin/virtlockd || { stop; exit 0; }

exec /usr/sbin/virtlockd -d
