# vnstat.conf
description "vnstat - lightweight network traffic monitor"

start on filesystem and net-device-up IFACE!=lo
stop on runlevel [016] or deconfiguring-networking

# TODO: readiness
respawn

setuid vnstat
setgid vnstat

pre-start script
    test -x /usr/sbin/vnstatd || { stop; exit 0; }
end script

exec /usr/sbin/vnstatd -n
