# vsftpd.conf
description "vsftpd - Very secure FTP server"

start on filesystem and static-network-up
stop on runlevel [016]

respawn

env CHROOT=/var/run/vsftpd/empty
env LOGFILE=/var/log/vsftpd.log

pre-start script
    [ -x /usr/sbin/vsftpd ] || { stop; exit 0; }
    [ -d $CHROOT ] || mkdir -p $CHROOT
    [ -f $LOGFILE ] || install -o root -g adm -m 640 $LOGFILE
end script

exec /usr/sbin/vsftpd
