# wicd.conf
description "wicd - network connection manager"

start on started dbus
stop on stopping dbus

# No readiness
respawn

pre-start script
    [ -x /usr/sbin/wicd ] || { stop; exit 0; }
    [ -d /var/run/wicd ] || mkdir -p /var/run/wicd
end script

exec /usr/sbin/wicd -f

post-stop exec rm -f /var/run/wicd/wicd.pid