# xl2tpd.conf
description "xl2tpd - layer 2 tunelling protocol daemon"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    [ -x /usr/sbin/xl2tpd ] || { stop; exit 0; }
    [ -d /var/run/xl2tpd ] || mkdir -p /var/run/xl2tpd
end script

exec /usr/sbin/xl2tpd
