# xl2tpd.conf
description "xl2tpd - layer 2 tunelling protocol daemon"

start on filesystem and started networking
stop on runlevel [016]

expect fork
respawn

pre-start script
    [ -x /usr/sbin/xl2tpd ] || { stop; exit 0; }
end script

script
    [ -s /etc/default/xl2tpd ] && . /etc/default/xl2tpd 
    [ -d ${XL2TPD_RUN_DIR:-/var/run/xl2tpd} ] || mkdir -p ${XL2TPD_RUN_DIR:-/var/run/xl2tpd}
    exec /usr/sbin/xl2tpd $DAEMON_OPTS
end script
