# xrdp-prep.conf
description "xrdp-prep - Preparation for X RDP server"

start on starting xrdp or starting xrdp-sesman

task

env RSAKEYS=/etc/xrdp/rsakeys.ini

script
	# PID file directory
	install -o xrdp -g xrdp -d /run/xrdp

	# Check for rsa key 
	test -s $RSAKEYS || {
		umask 077
		xrdp-keygen xrdp $RSAKEYS
		chown xrdp:xrdp $RSAKEYS
		test -f $RSAKEYS || { echo "could not create $RSAKEYS"; exit 1; }
	}
end script
