# xrdp-sesman.conf
description "xrdp-sesman - X RDP Session Manager"

start on started xrdp
stop on stopping xrdp

respawn
exec /usr/sbin/xrdp-sesman -n
