# xrdp.conf
description "xrdp - a Remote Desktop Protocol (RDP) server"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

setuid xrdp
setgid xrdp

respawn
exec /usr/sbin/xrdp -nodaemon
