# xringd.conf
description "xringd - Phone line monitor"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn
reload signal 0

pre-start script
    test -x /usr/sbin/xringd || { stop; exit 0; }
end script

exec /usr/sbin/xringd
