# yate.conf
description "YATE - Yet Another Telephony Engine"

start on filesystem and net-device-up IFACE=lo
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

setuid yate
setgid yate

# Change in /etc/default/yate to be able to start service
env NO_START=1

pre-start script
    test -x /usr/bin/yate || { stop; exit 0; }
    if [ -r /etc/default/yate ]; then
        . /etc/default/yate
    fi
    [ "$NO_START" != "1" ] || { stop; exit 0; }
end script

exec /usr/bin/yate -d -rs -l /var/log/yate/yate.log -Dfu
