# zabbix-agent.conf
description "zabbix-agent - network monitoring solution (agent)"

start on runlevel [2345]
stop on runlevel [016] or unmounting-filesystem

pre-start exec install -o zabbix -g zabbix -d /run/zabbix

respawn
expect daemon
exec /usr/sbin/zabbix_agentd
