# zephyrd.conf
description "zephyrd - Zephyr IM server"

start on filesystem and static-network-up
stop on deconfiguring-networking or runlevel [016]

expect fork
respawn

pre-start script
    test -x /usr/sbin/zephyrd || { stop; exit 0; }
end script

exec /usr/sbin/zephyrd
