# base/nldev/devtriggerall.conf
description "devtriggerall: trigger add events for coldplug devices"

start on started nldev and not-container
emits coldplug-trigger

task
exec devtriggerall
post-start exec initctl emit -n coldplug-trigger
