# base/nldev/nldev.conf
description "nldev: netlink device manager"

start on virtual-filesystems
stop on unmounting-filesystem
emits *-device-*

respawn
exec nldev -k -r /sbin/upstart-mdev-bridge
