# base/nldev/nltrigger.conf
start on started nldev and not-container
emits coldplug-trigger

task
exec nltrigger add all
post-start exec initctl emit -n coldplug-trigger
