# print/cups-browsed.conf
description "cups-browsed - Make remote CUPS printers available locally with Bonjour"

apparmor load /etc/apparmor.d/usr.sbin.cups-browsed

start on started avahi-daemon
stop on (runlevel [016] or unmounting-filesystem or
         stopping avahi-daemon)

respawn
respawn limit 3 240

exec cups-browsed
