# rpc/gssd-mounting.conf
# gssd-mounting

description "Block the mounting event for NFS4 filesytems until gssd is running"
author "Steve Langasek <steve.langasek@canonical.com>"

instance $MOUNTPOINT

start on mounting TYPE=nfs* OPTIONS=*sec*krb5*
stop on started gssd or stopped gssd
task

# This is required so that the task is still considered
# successful when it gets killed
normal exit TERM

script

  DEFAULTFILE=/etc/default/nfs-common

  if [ -f "$DEFAULTFILE" ]; then
    . "$DEFAULTFILE"
  fi

  if [ -f /etc/fstab ]; then
    exec 9<&0 </etc/fstab

    while read DEV MTPT FSTYPE OPTS REST
    do
      case "$OPTS" in
        sec=krb5|*,sec=krb5|sec=krb5,*|*,sec=krb5,*|sec=krb5i|*,sec=krb5i|sec=krb5i,*|*,sec=krb5i,*|sec=krb5p|*,sec=krb5p|sec=krb5p,*|*,sec=krb5p,*)
          AUTO_NEED_GSSD=yes
          ;;
      esac
    done

    exec 0<&9 9<&-
  fi

  case "$NEED_GSSD" in
    yes|no)
      ;;
    *)
      NEED_GSSD=$AUTO_NEED_GSSD
      ;;
  esac
  [ "x$NEED_GSSD" = xyes ] || { stop; exit 0; }

  status gssd | grep -q "start/running" && exit 0

  # If it's already starting we'll get killed by the impending 'stop on
  # started gssd' 
  # If it wasn't already starting, we'll either get killed by the stop
  # on started or stopped.
  # So, its safe to sleep forever here and rely on upstart to kill us,

  while sleep 3600; do :; done
end script
