# rpc/idmapd.conf
# idmapd - NFSv4 id <-> name mapper

# rpc.idmapd is the NFSv4 ID <-> name mapping daemon.  It provides
# functionality to the NFSv4 kernel client and server, to which it
# communicates via upcalls, by translating user and group IDs to names, and
# vice versa.

description	"NFSv4 id <-> name mapper"
author		"Steve Langasek <steve.langasek@canonical.com>"

start on local-filesystems
stop on unmounted-remote-filesystems

expect fork
respawn

env PIPEFS_MOUNTPOINT=/run/rpc_pipefs

pre-start script
	do_modprobe() {
		modprobe -q "$1" || true
	}

	do_modprobe nfs
	do_modprobe nfsd

	do_modprobe sunrpc

	if ! mountpoint -q "$PIPEFS_MOUNTPOINT"
	then
		mkdir -p "$PIPEFS_MOUNTPOINT"
		mount -t rpc_pipefs rpc_pipefs "$PIPEFS_MOUNTPOINT" || true
	fi
end script

exec rpc.idmapd

post-stop script
	if mountpoint -q "$PIPEFS_MOUNTPOINT"
	then
		# ignore any failures caused by the filesystem still
		# being in use
		umount "$PIPEFS_MOUNTPOINT" || true
	fi
end script

