# rpc/start-ypbind.conf
description	"block processes that need ypbind"
author		"Steve Langasek <steve.langasek@canonical.com>"

instance $JOB
env JOB= 

start on (starting autofs or starting am-utils or starting gdm
          or starting lightdm or starting atd)

task

exec start wait-for-state WAITER="$JOB" WAIT_FOR=ypbind WAIT_STATE=running
