# rpc/yppasswdd.conf
description	"NIS password update daemon"
author		"Steve Langasek <steve.langasek@canonical.com>"

start on started ypserv
stop on stopping ypserv

respawn
expect daemon

pre-start script
	NISSERVER=false

	[ -f /usr/sbin/ypbind ] && [ -f /etc/defaultdomain ] || { stop; exit 0; }

	[ -f /etc/default/nis ] && . /etc/default/nis

	[ "$NISSERVER" = master ] || { stop; exit 0; }
end script

script
	YPCHANGEOK=chsh
	YPPWDDIR=/etc

	[ -f /etc/default/nis ] && . /etc/default/nis

	E=
	if [ -n "$YPCHANGEOK" ]; then
		OIFS="$IFS"; IFS="$IFS,"
		for i in $YPCHANGEOK
		do
			case $i in
				chsh|chfn)
					E="$E -e $i"
					;;
			esac
		done
		IFS="$OIFS"
	fi
	if [ -n "$YPPWDDIR" ]; then
		YPPWDDIRARGS="-D $YPPWDDIR"
	fi
	exec rpc.yppasswdd $YPPWDDIRARGS $E $YPPASSWDDARGS
end script
