# rpc/ypserv.conf
description	"NIS server"
author		"Steve Langasek <steve.langasek@canonical.com>"

# Patterned after the statd.conf job in nfs-common.
start on (started portmap ON_BOOT=
          or (started portmap ON_BOOT=y
              and ((filesystem and static-network-up) or failsafe-boot)))
stop on stopping portmap or runlevel [!2345]

expect daemon
respawn

pre-start script
	NISSERVER=false

	[ -f /usr/sbin/ypbind ] && [ -f /etc/defaultdomain ] || { stop; exit 0; }

	[ -f /etc/default/nis ] && . /etc/default/nis

	[ "$NISSERVER" = false ] && { stop; exit 0; }

	oname=$(domainname)
	nname=$(cat /etc/defaultdomain)
	if [ "$oname" != "$nname" ]; then
		echo "Setting NIS domainname to: $nname"
		domainname "$nname"
	fi
end script

script
	[ -f /etc/default/nis ] && . /etc/default/nis
	exec ypserv $YPSERVARGS
end script

post-start script
	[ -f /etc/default/nis ] && . /etc/default/nis
	if [ "$NISSERVER" = slave ] && [ -n "$NISMASTER" ]; then
		/usr/lib/yp/ypinit -s "$NISMASTER"
	fi
end script
