# rpc/ypxfrd.conf
description	"NIS map transfer server"
author		"Steve Langasek <steve.langasek@canonical.com>"

start on started ypserv
stop on stopping ypserv

respawn
expect daemon

pre-start script
	NISSERVER=false

	[ -f /usr/sbin/ypbind ] && [ -f /etc/defaultdomain ] || { stop; exit 0; }

	[ -f /etc/default/nis ] && . /etc/default/nis

	[ "$NISSERVER" = master ] || { stop; exit 0; }
end script

script
	[ -f /etc/default/nis ] && . /etc/default/nis

	exec rpc.ypxfrd $YPXFRDARGS
end script
