# message-servers/ejabberd.conf
description "ejabberd - a distributed fault-tolerant Jabber/XMPP server"

start on runlevel [2345]
stop on runlevel [016] or deconfiguring-networking or unmounting-filesystem

chdir /var/lib/ejabberd

pre-start exec install -o ejabberd -g ejabberd -m 755 -d /run/ejabberd

respawn
exec start-stop-daemon -S -q -c ejabberd -x /usr/sbin/ejabberd -- -noshell
