# message-servers/quasselcore.conf
# quasselcore - distributed IRC client using a central core component
#
# This is the core component of Quassel. A modern,
# cross-platform, distributed IRC client, meaning that one
# (or multiple) client(s) can attach to and detach from this
# central core. It's much like the popular combination
# of screen and a text-based IRC client such as WeeChat.

description	"distributed IRC client using a central core component"
author		"Michael Marley <michael@michaelmarley.com>"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [!2345]

setuid quasselcore
setgid quassel

chdir /var/lib/quassel
env HOME=/var/lib/quassel

respawn

script
	DATADIR="/var/lib/quassel"
	LOGFILE="/var/log/quassel/core.log"
	LOGLEVEL="Info"
	PORT="4242"

	test -f /etc/default/quasselcore && . /etc/default/quasselcore

	exec quasselcore --configdir=$DATADIR --logfile=$LOGFILE --loglevel=$LOGLEVEL --port=$PORT
end script
