# INSTALL.minimal/mountall/mounted/container-detect.conf
description "container-detect"

start on mounted MOUNTPOINT=/run

env container
env LIBVIRT_LXC_UUID

emits container
emits not-container

pre-start script
    # Detect old-style libvirt
    if [ -z "$container" ] && [ -n "$LIBVIRT_LXC_UUID" ]; then
        container=lxc-libvirt
    fi

    # Detect OpenVZ containers
    if [ -z "$container" ] && [ -d /proc/vz ] && [ ! -d /proc/bc ]; then
        container=openvz
    fi

    # Detect Vserver containers
    if [ -z "$container" ]; then
        VXID="$(cat /proc/self/status | grep ^VxID | cut -f2)" || true
        [ "${VXID:-0}" -le 1 ] || container=vserver
    fi

    if [ -n "$container" ]; then
        echo "$container" > /run/container_type || true
        initctl emit -n container CONTAINER=$container
    else
        rm -f /run/container_type
        initctl emit -n not-container
        stop
    fi
end script
