# INSTALL.minimal/mountall/mounted/mounted-debugfs.conf
description "mounted-debugfs: restrict debugfs permissions"

start on mounted MOUNTPOINT=/sys/kernel/debug TYPE=debugfs
env MOUNTPOINT=/sys/kernel/debug

task

script
    chmod 0700 "${MOUNTPOINT}" || true
end script
