# INSTALL.minimal/mountall/mounted/mounted-tmp.conf
description "mounted-tmp: clean /tmp if not on a tmpfs"

start on mounted MOUNTPOINT=/tmp TYPE!=tmpfs
env MOUNTPOINT=/tmp

task

script
	cd "${MOUNTPOINT}" || exit 1
	rm -f .X*-lock

	[ -r /etc/default/rcS ] && . /etc/default/rcS
	case "${TMPTIME}" in
	  -*|infinite|infinity)
	    exit 0
	    ;;
	esac

	if [ "${TMPTIME}" != "0" -a -n "${TMPTIME}" ]; then
		DEXPR="-mtime +${TMPTIME} -ctime +${TMPTIME}"
		TEXPR="$DEXPR -atime +${TMPTIME}"
	fi

	EXCEPT='! -name . ! ( -path "./...security*" -uid 0)'
	paths="./lost+found ./quota.user ./aquota.user ./quota.group ./aquota.group ./.journal ./.clean"
	for p in $paths; do
		EXCEPT="$EXCEPT ! ( -path ${p} -uid 0 )"
	done

	# Remove all old files, then all empty directories
	find . -depth -xdev $TEXPR $EXCEPT ! -type d -delete
	find . -depth -xdev $DEXPR $EXCEPT -type d -empty -delete

	# Check if we have enough space in /tmp, and if not, mount a tmpfs there
	avail=`df -kP /tmp | awk 'NR==2 { print $4 }'`
	if [ "$avail" -lt 1000 ]; then
		mount -t tmpfs -o size=1048576,mode=1777 overflow /tmp
	fi
end script
