# INSTALL.minimal/mountall/mounted/mounted-var.conf
description "mounted-var: add run and lock links to separate /var filesystem"

start on mounted MOUNTPOINT=/var

task

script
    run_migrate () {
	OLD="$1"
	RUN="$2"

	if [ -L "$OLD" ]; then
		[ "$(readlink "$OLD")" != "$RUN" ] || return 0
		rm -f "$OLD"
	elif [ -d "$OLD" ]; then
		rm -rf "$OLD" || :
	fi

	ln -fs "$RUN" "$OLD"
	[ ! -x /sbin/restorecon ] || /sbin/restorecon "$OLD"
    }

    run_migrate /var/run /run
    run_migrate /var/lock /run/lock
end script
