# INSTALL.minimal/mountall/mounted/mounted-dev.conf
description "mounted-dev: create /dev/shm symlink"

start on mounted MOUNTPOINT=/dev
env MOUNTPOINT=/dev

task

script
    [ -e /dev/shm ] || ln -s /run/shm /dev/shm
end script
