# fs/mountall/mounted/mounted-proc.conf
description	"mounted-proc: restrict access to /proc/slabinfo"

start on mounted MOUNTPOINT=/proc TYPE=proc
env MOUNTPOINT=/proc

task

script
    test ! -e "${MOUNTPOINT}"/slabinfo  || chmod 0400 "${MOUNTPOINT}"/slabinfo
end script
