# fs/mountall/mounted/mounted-run.conf
description "mounted-run: populate and link to /run filesystem"

start on mounted MOUNTPOINT=/run TYPE=tmpfs

pre-start script
	: > "/run/utmp"
	chmod 664 "/run/utmp"
	chgrp utmp "/run/utmp"
	mkdir -p /run/sendsigs.omit.d
end script
