# fs/mountall/mounted/mounted-var.conf
description "mounted-var: add run and lock links to separate /var filesystem"

start on mounted MOUNTPOINT=/var

task

script
    run_migrate () {
	OLD="$1"
	RUN="$2"

	if [ -L "$OLD" ]; then
	    if [ "$(readlink "$OLD")" != "$RUN" ]; then
		# Remove any old (relative?) symlinks
		rm -f "$OLD"
	    fi
	elif [ -d "$OLD" ]; then
	    # Remove old directories that are explicitly ephemeral
	    rm -rf "$OLD" 2>/dev/null || true
	fi

	if [ ! -L "$OLD" ]; then
	    ln -fs "$RUN" "$OLD"
	    [ -x /sbin/restorecon ] && /sbin/restorecon "$OLD"
	fi

	return 0
    }

    run_migrate /var/run /run
    run_migrate /var/lock /run/lock
end script
