# dev/busybox-uevent.conf
description "busybox-uevent: run mdev on netlink events"

start on virtual-filesystems
stop on unmounting-filesystem
emits *-device-*

respawn
exec busybox uevent upstart-mdev-bridge
post-start exec busybox mdev -s
