# dev/coldplug-trigger.conf
description "coldplug-trigger: synchronization job for sysinit checkpoint"

# blocks sysinit-up event until we see the coldplug-trigger
start on sysinit-up and coldplug-trigger
