# dev/devtriggerall.conf
description "uevent-trigger: add events for coldplug devices"

start on not-container and started busybox-uevent
emits coldplug-trigger

task
exec find /sys/devices -name uevent -exec sh -c 'echo -n add > {}' \;
post-start exec initctl emit -n coldplug-trigger
