# dev/devtriggerall.conf
description "devtriggerall: trigger add events for coldplug devices"
# upstream: https://github.com/CameronNemo/linux-device-tools

start on started busybox-uevent and not-container
emits coldplug-trigger

task
exec devtriggerall
post-start exec initctl emit -n coldplug-trigger
