# dev/faketrigger.conf
description "faketrigger: emit coldplug event in containers"
start on container
emits coldplug-trigger
task
exec initctl emit -n coldplug-trigger
