# acpi/acpid.conf
description "acpid: framework for reacting to acpi events"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

exec acpid -c /etc/acpi/events -s /var/run/acpid.socket
