# misc/sympa/sympa-bounced.conf
description "sympa - Mailing List Manager (bounced)"

start on start-sympa
stop on stop-sympa

setuid sympa
setgid sympa

respawn
exec /usr/lib/sympa/bin/bounced.pl -F
