# misc/sympa/sympa-archived.conf
description "sympa - Mailing List Manager (archived)"

start on start-sympa
stop on stop-sympa

setuid sympa
setgid sympa

respawn
exec /usr/lib/sympa/bin/archived.pl -F
