# misc/sympa/sympa.conf
description "Sympa - modern mailing list manager"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

pre-start script
	install -o sympa -g sympa -m 755 -d /run/sympa
	exec initctl emit start-sympa
end script

post-stop exec initctl emit stop-sympa
