# misc/sympa/sympa-main.conf
description "sympa - Mailing List Manager (main daemon)"

start on start-sympa
stop on stop-sympa

setuid sympa
setgid sympa

respawn
exec /usr/lib/sympa/bin/sympa.pl --foreground
