# misc/sympa/sympa-task_manager.conf
description "sympa - Mailing List Manager (task_manager)"

start on start-sympa
stop on stop-sympa

setuid sympa
setgid sympa

respawn
exec /usr/lib/sympa/bin/task_manager.pl -F
