# base/udev/udev-fallback-graphics.conf
# fallback-graphics - take actions to initiate fallback graphics
#
# if we have cold plugged everything and not yet seen a valid graphics
# device,
# shove a fallback framebuffer into the mix and hope for the best.

description	"load fallback graphics devices"

# We only want this job to happen once per boot,
# we can use a pre-start stanza, or
# make it a task and start on 'startup and ...'.
start on (graphics-device-added PRIMARY_DEVICE_FOR_DISPLAY=1
          or drm-device-added PRIMARY_DEVICE_FOR_DISPLAY=1
          or started udevtrigger)

pre-start script
    if [ "$PRIMARY_DEVICE_FOR_DISPLAY" = "" ]; then
        modprobe -q -b vesafb || true
    fi
end script
