# base/udev/udevtrigger.conf
# udevtrigger - cold plug devices
#
# By the time udevd starts, we've already missed all of the events for
# the devices populated in /sys.  This task causes the kernel to resend
# them.

description	"cold plug devices"

start on started udev

pre-start script
    # no cold plug in containers
    test ! -f /run/container_type || exit 0
    # trigger the add events, then wait for them to settle
    udevadm trigger --action=add
    exec udevadm settle
end script
