# sessions/gnome-shell.conf
# BUG: natural scrolling not working for panel, but is working for activities menu
# Randomly starts working after some time...

start on graphical-session-start SESSION=gnome and gnome-session-up
stop on graphical-session-end and gnome-session-down

# Compatibility with old Qt programs
env GNOME_DESKTOP_SESSION_ID=this-is-deprecated
export GNOME_DESKTOP_SESSION_ID

# Avoid pollution from start events
exec env -u UPSTART_JOB -u UPSTART_INSTANCE -u UPSTART_EVENTS -u SESSION \
	/usr/bin/gnome-shell
