# sessions/gnome-shell-session.conf
start on graphical-session-start SESSION=gnome and started dbus
stop on graphical-session-end

emits gnome-session-{start,up,end,down}

pre-start exec /sbin/initctl emit gnome-session-start 
post-start exec /sbin/initctl emit -n gnome-session-up

post-stop script
	/sbin/initctl emit gnome-session-end
	/sbin/initctl emit -n gnome-session-down
end script
