# sessions/gnome-session-bin.conf
# BUG: window state is not saved

description "Gnome Session itself"
author "Stéphane Graber <stgraber@ubuntu.com>"
# with mods by Cameron Norman <camerontnorman@gmail.com>

start on gnome-session-start
stop on gnome-session-end

emits session-manager-up
emits session-end

env GNOME_DESKTOP_SESSION_ID=this-is-deprecated
export GNOME_DESKTOP_SESSION_ID

exec env -u UPSTART_EVENTS -u UPSTART_INSTANCE -u UPSTART_JOB gnome-session --session=upstart

# Set SESSION_MANAGER (hack)
post-start script
	pid=$(initctl status | grep start/ | cut -d" " -f4)
	base="$(hostname):@/tmp/.ICE-unix/$pid"
	initctl set-env --global SESSION_MANAGER="local/$base,unix/$base"
	initctl emit -n session-manager-up
end script

# GNOME session exits to signal when to log out
post-stop exec initctl emit -n session-end
