# sessions/jackd.conf
start on startup
stop on session-end

exec jackd -R -P 70 -d alsa
