# sessions/light-locker.conf
description "Lock screen that uses LightDM"

start on graphical-session-start
stop on graphical-session-end

respawn
exec /usr/bin/light-locker
