# sessions/nm-applet.conf
start on graphical-session-start SESSION!=gnome
stop on graphical-session-end

exec /usr/bin/nm-applet
