# sessions/polkit-gnome-authentication-agent-1.conf
start on graphical-session-start SESSION!=gnome
stop on graphical-session-end

exec /usr/lib/policykit-1-gnome/polkit-gnome-authentication-agent-1
