# sessions/print-applet.conf
start on gnome-session-start or starting gnome-session-bin
stop on gnome-session-end

respawn
exec /usr/share/system-config-printer/applet.py
