# sessions/pulseaudio.conf
# Debian: this does not account for
# start-pulseaudio-kde's special module loading

# You must disable auto spawning of PulseAudio in:
# /etc/pulseaudio/client.conf for this to work correctly

# Auto start entries (/etc/xdg/autostart, ~/.config/autostart)
# also need to be disabled or removed.

description "PulseAudio - sound server"

start on startup
stop on session-end

# Just in case PA was somehow started under non-upstart
pre-start script
	pkill pulseaudio || true
end script

respawn

exec pulseaudio --start --log-target=stderr --daemonize=false \
		--load "module-x11-publish display=$DISPLAY" \
		--load "module-x11-cork-request display=$DISPLAY"
