# sessions/xsession-init.conf
description "XSession handling"
author "Stéphane Graber <stgraber@ubuntu.com>"

start on startup
stop on session-end

emits graphical-session-{start,end}
emits desktop-session-{start,end}

pre-start script
	test -n "$DESKTOP_SESSION" || { /sbin/stop; exit 0; }
	initctl emit graphical-session-start SESSION=$DESKTOP_SESSION SESSIONTYPE=$SESSIONTYPE || true
	initctl emit desktop-start DESKTOP_SESSION=$DESKTOP_SESSION || true
end script

post-stop script
	# Don't wait to avoid blocking the session from ending.
	initctl emit --no-wait desktop-end DESKTOP_SESSION=$DESKTOP_SESSION || true
	graphical-session-end SESSION=$DESKTOP_SESSION SESSIONTYPE=$SESSIONTYPE || true
	exec dbus-send --type=method_call --address=$UPSTART_SESSION /com/ubuntu/Upstart com.ubuntu.Upstart0_6.EndSession
end script
