# task-managers/atd.conf
description "atd: Deferred execution scheduler"

start on filesystem
stop on unmounting-filesystem

expect fork
respawn

exec atd
