# task-managers/monit.conf
description "monit - service and resource monitoring daemon"

start on filesystem and started networking
stop on runlevel [016] or unmounting-filesystem

env CONFIG=/etc/monit/monitrc

respawn
script
	[ -s /etc/default/monit ] && . /etc/default/monit
	exec /usr/bin/monit -Ic $CONFIG $MONIT_OPTS
end script
