# task-managers/svscan.conf
# svscan - daemontools
#
# This service starts daemontools from the point the system is
# started until it is shut down again.

start on filesystem
stop on runlevel [06]

respawn
exec /usr/bin/svscanboot
