# task-managers/xinetd.conf
description "xinetd - daemon superserver"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

pre-start script
    test -x /usr/sbin/xinetd || { stop; exit 0; }
end script

exec /usr/sbin/xinetd -stayalive -inetd_compat -inetd_ipv6
