# backup/bacula/bacula-director.conf
description "bacula-director - supervises bacula file operations"

# TODO: start after mysql and postgresql have started (somehow)
start on (filesystem
          or started bacula-fd)
stop on (runlevel [016] or unmounting-filesystem
         or stopping bacula-fd)

expect fork
respawn

env RUNDIR=/var/run/bacula
env CONF=/etc/bacula/bacula-sd.conf

pre-start script
    [ -x /usr/sbin/bacula-sd ] || { stop; exit 0; }
    [ -d $RUNDIR ] || install -d -o bacula -g daemon -m 0755 $RUNDIR
end script

exec /usr/sbin/bacula-sd -c $CONF -u bacula -g tape
