# backup/bacula/bacula-fd.conf
description "bacula-fd - Bacula File Daemon"

start on (filesystem
          or starting bacula-director
          or started bacula-sd)
stop on runlevel [016] or unmounting-filesystem or stopping bacula-sd

expect fork
respawn

env RUNDIR=/var/run/bacula
env CONF=/etc/bacula/bacula-fd.conf

pre-start script
    [ -x /usr/sbin/bacula-fd ] || { stop; exit 0; }
    [ -d $RUNDIR ] || install -d -o bacula -g daemon -m 0755 $RUNDIR
end script

exec /usr/sbin/bacula-fd -c $CONF -u bacula
