# ltsp/ltsp-cluster-lbagent.conf
# ltsp-cluster - LoadBalancer agent
#

description "ltsp-cluster LoadBalancer agent"
author      "Stephane Graber <stgraber@ubuntu.com>"

start on (started ltsp-cluster-accountmanager)
stop on runlevel[016]

expect daemon
respawn

exec ltsp-cluster-lbagent --daemon --logfile=/var/log/ltsp-cluster-lbagent.log /etc/ltsp/lbaconfig.xml
