# ltsp/ltsp-cluster-lbserver.conf
# ltsp-cluster - Lbserver
#

description "ltsp-cluster-lbserver"
author      "Marc Gariépy <mgariepy@ubuntu.com>"

start on net-device-up and local-filesystems
stop on runlevel[016]

respawn

exec /usr/sbin/ltsp-cluster-lbserver --logfile=/var/log/ltsp-cluster-lbserver.log /etc/ltsp/lbsconfig.xml

