# ltsp/ltsp-cluster-accountmanager.conf
# ltsp-cluster - Accountmanager
#

description "ltsp-cluster-accountmanager"
author      "Marc Gariépy <mgariepy@ubuntu.com>"

start on net-device-up and local-filesystems and (started dbus)
stop on runlevel[016]

respawn

exec /usr/bin/ltsp-cluster-accountmanager /etc/ltsp/accountmanager.conf
