# ltsp/ltsp-cluster-agent.conf
description "LTSP Agent"
author      "Stephane Graber <stgraber@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec ltsp-agent
