# monitoring/collectd.conf
description "collectd - manage the statistics collection daemon"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn

env CONFIG=/etc/collectd/collectd.conf

pre-start script
    [ -x /usr/sbin/collectd ] || { stop; exit 0; }
    [ -s /etc/default/collectd ] && . /etc/default/collectd
    [ "$DISABLE" = "0" ] || { stop; exit 0; }
    # Choke if config file is invalid or non-existent
    [ -s $CONFIG ] || { stop; exit 0; }
    collectd -t -C $CONFIG || stop; exit 0; }
end script

exec collectd -C $CONFIG
