# monitoring/elasticsearch.conf
description "elasticsearch"

start on started nss-lookup
stop on stopping nss-lookup

env ES_HOME=/usr/share/elasticsearch
env ES_PATH_CONF=/etc/elasticsearch

chdir /usr/share/elasticsearch
setuid elasticsearch
setgid elasticsearch

# Specifies the maximum file descriptor number that can be opened by this process
limit nofile 65536
# Specifies the maximum number of processes
limit nproc 4096
# Specifies the maximum size of virtual memory
limit as unlimited
# Specifies the maximum file size
limit fsize unlimited

kill signal TERM
kill timeout 0
normal exit 143
exec /usr/share/elasticsearch/bin/elasticsearch -p /run/elasticsearch/elasticsearch.pid --quiet
