# monitoring/elasticsearch-rundir.conf
start on starting elasticsearch
task
exec install -o elasticsearch -g elasticsearch -m 0755 -d /run/elasticsearch
