# monitoring/grafana.conf
description "grafana: monitoring and metric analytics"

start on started sysinit
stop on stopping sysinit

setuid grafana
setgid grafana
chdir /usr/share/grafana
limit nofile 10000
umask 027

env CONF_FILE=/etc/grafana/grafana.ini
env LOG_DIR=/var/log/grafana
env DATA_DIR=/var/lib/grafana
env PLUGINS_DIR=/var/lib/grafana/plugins
env PROVISIONING_CFG_DIR=/etc/grafana/provisioning

kill timeout 20
respawn
exec /usr/sbin/grafana-server --config=${CONF_FILE}                                   \
                              --pidfile=/run/grafana-server.pid                       \
                              cfg:default.paths.logs=${LOG_DIR}                       \
                              cfg:default.paths.data=${DATA_DIR}                      \
                              cfg:default.paths.plugins=${PLUGINS_DIR}                \
                              cfg:default.paths.provisioning=${PROVISIONING_CFG_DIR}  
