# monitoring/grafana-pre.conf
description "grafana-pre: setup and DB wait"

start on starting grafana

task

pre-start script
	install -m 0750 -o grafana -g grafana -d /run/grafana

	# wait for database(s) to come online
	for svc in "postgresql mariadb mysql"; do
		stat="$(status ${svc})"

		# skip if the job is not installed
		test "$?" = 0 || continue

		while true; do
			# continue if the job's goal is `stop` or state is `running`
			case "$stat" in
				*stop/*)
				*start/running*)
					break
					;;
				*)
					sleep 0.666
					;;
			esac
		done
	done
end script
