# monitoring/icinga.conf
description "icinga: monitoring and management system"

start on started nss-lookup
stop on stopping nss-lookup

nice 5
env CONF="/etc/icinga/icinga.cfg"

pre-start script
	test -x /usr/sbin/icinga || { stop; exit 0; }
	install -m 0750 -o nagios -g nagios -d /var/run/icinga
	nagios_pipe=$(grep "^command_file" "$CONF" | sed 's@=@ @')
	test ! -e "$nagios_pipe" || rm -f "$nagios_pipe"
end script

exec /usr/sbin/icinga "$CONF"
