# monitoring/kibana.conf
description "kibana: analytics and search dashboard for Elasticsearch"

start on started sysinit
stop on stopping sysinit

setuid kibana
setgid kibana

respawn
exec /usr/share/kibana/bin/kibana -c /etc/kibana/kibana.yml
