# monitoring/mtail.conf
description "mtail: extract monitoring data from logs"

start on started sysinit
stop on stopping sysinit

setuid mtail

pre-start script
	[ ! -s /etc/default/mtail ] || . /etc/default/mtail
	if [ -z "$ENABLED" -o "$ENABLED" = "0" -o "$ENABLED" = "no" ] ; then
		echo "disabled, see /etc/default/mtail"
		stop; exit 0
	fi
end script

script
        [ ! -s /etc/default/mtail ] || . /etc/default/mtail

	ARGS="-progs /etc/mtail -logtostderr"

	[ -n "$LOGS" ] && ARGS="${ARGS} -logs $LOGS"
	[ -n "$PORT" ] && ARGS="${ARGS} -port $PORT"
	[ -n "$COLLECTD_SOCKETPATH" ] && \
	    ARGS="${ARGS} -collectd_socketpath $COLLECTD_SOCKETPATH"
	[ -n "$GRAPHITE_HOSTPORT" ] && \
	    ARGS="${ARGS} -graphite_host_port $GRAPHITE_HOSTPORT"
	[ -n "$STATSD_HOSTPORT" ] && ARGS="${ARGS} -statsd_hostport $STATSD_HOSTPORT"
	[ -n "$METRIC_PUSH_INTERVAL" ] && \
	    ARGS="${ARGS} -metric_push_interval_seconds $METRIC_PUSH_INTERVAL"
	[ -n "$EXTRA_ARGS" ] && ARGS="${ARGS} $EXTRA_ARGS"

	exec /usr/bin/mtail "$ARGS"
end script
