# monitoring/nagios3.conf
description "nagios3: monitoring and management system"

start on started nss-lookup
stop on stopping nss-lookup

nice 5
env CONF="/etc/nagios3/nagios.cfg"

pre-start script
	test -x /usr/sbin/nagios3 || { stop; exit 0; }
	install -m 0750 -o nagios -g nagios -d /var/run/nagios3
	nagios_pipe=$(grep "^command_file" "$CONF" | sed 's@=@ @')
	test ! -e "$nagios_pipe" || rm -f "$nagios_pipe"
end script

exec /usr/sbin/nagios3 "$CONF"
