# monitoring/prometheus.conf
description "prometheus: monitoring system and time series database"

start on started sysinit
stop on stopping sysinit

setuid prometheus

kill timeout 20
respawn
exec /usr/bin/prometheus
