# monitoring/sensu.conf
description "sensu: full stack monitoring"

start on started sysinit
stop on stopping sysinit

pre-start script
	install -o sensu -g sensu -m 0755 -d /run/sensu
	exec >/dev/null 2>&1
	start sensu-service SVC=api || :
	start sensu-service SVC=client || :
	start sensu-service SVC=server || :
end script
