# monitoring/uchiwa.conf
description "uchiwa: dashboard for sensu"

start on started sysinit
stop on stopping sysinit

setuid uchiwa
chdir /opt/uchiwa/src

exec /opt/uchiwa/bin/uchiwa -c /etc/sensu/uchiwa.json -d /etc/sensu/dashboard.d -p /opt/uchiwa/src/public
