# swift/swift-object-expirer.conf
description "SWIFT Object Expirer"
author "Will Kelly <william.kelly@rackspace.com>"

start on runlevel [2345]
stop on runlevel [016]

pre-start script
  if [ -f "/etc/swift/object-expirer.conf" ]; then
    exec /usr/bin/swift-init object-expirer start
  else
    exit 1
  fi
end script

post-stop exec /usr/bin/swift-init object-expirer stop
