# database/trove/trove-conductor.conf
description "Trove conductor server"
author "Chuck Short <zulcss@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start --chuid trove:trove \
            --chdir /var/lib/trove --name trove-conductor \
            --exec /usr/bin/trove-conductor -- --config-file=/etc/trove/trove-conductor.conf --log-dir=/var/log/trove --logfile=conductor.log
