# database/trove/trove-taskmanager.conf
description "Trove API server"
author "Thomas Goirand <zigo@debian.org>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start --chuid trove:trove \
            --chdir /var/lib/trove --name trove-api \
            --exec /usr/bin/trove-taskmanager -- --config-file=/etc/trove/trove-taskmanager.conf --log-dir=/var/log/trove --logfile=taskmanager.log
