# broken/4store.conf
description "4store - database server daemon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

env USER=fourstore

pre-start script
    test -x /usr/bin/4s-backend -a -x /usr/bin/4s-httpd || { stop; exit 0; }
    install -m 0755 -o $USER -g $USER -d /var/lib/4store /var/log/4store
    exec initctl emit start-fourstore
end script
