# broken/dbeacon.conf
description "dbeacon - Multicast Beacon"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and static-network-up
stop on runlevel [016] or deconfiguring-networking

pre-start script
    test -x /usr/bin/dbeacon || { stop; exit 0; }
    SUCCESS=true
    for i in /etc/dbeacon/*.conf; do
        i=${i#/etc/dbeacon/}
        i=${i%.conf}
        start dbeacon-instance INSTANCE=$i || SUCCESS=false
    done
    test "$SUCCESS" = "true" || exit 1
end script

# Instances stop when this stops
