# broken/fail2ban.conf
# This job needs a modification in the file /usr/bin/fail2ban-client
# Specifically, it needs to pass the option "-f" to the server, not "-b"
# Search for args.append("-b") in the file to find the place to change this

description "fail2ban - ban hosts that cause multiple authentication errors"

start on filesystem and started networking
stop on deconfiguring-networking

expect fork
respawn

env RUNDIR=/var/run/fail2ban

pre-start script
    [ -d $RUNDIR ] || mkdir -p $RUNDIR
end script

exec /usr/bin/fail2ban-client -x -f start

pre-stop exec /usr/bin/fail2ban-client stop
