# broken/ipsec.conf
description "ipsec - Strongswan IPsec services"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem and started networking
stop on runlevel [016]

pre-start script
    mkdir -p /var/lock/subsys
    
    exec ipsec start
end script

pre-stop exec ipsec stop
