# broken/lirc-lircd.conf
description "lirc-lircd - main remote control daemon of LIRC"

start on starting lirc
stop on stopping lirc

expect fork
respawn

pre-start script
    [ -x /usr/sbin/lircd ] || { stop; exit 0; }

    if [ ! -f /etc/lirc/lircd.conf ] || \
       grep -q "^#UNCONFIGURED" /etc/lirc/lircd.conf
    then
        echo "Not configured; not starting."
        stop; exit 0
    fi

    [ -d /var/run/lirc ] || mkdir -p /var/run/lirc
end script

script
    [ -s /etc/lirc/hardware.conf ] && . /etc/lirc/hardware.conf
    exec /usr/sbin/lircd --driver="$DRIVER" --device="$DEVICE"
end script

post-start exec ln -fs /var/run/lirc/lircd /dev/lircd
