# broken/lirc-lircmd.conf
description "lirc-lircmd - LIRC remote control mouse daemon"

# This start on ensures that anything that wants to start after all of
# these daemons can do so by simply using the event "started lirc".
start on starting lirc and started lirc-lircd
stop on stopping lirc

expect fork
respawn

pre-start script
    [ -x /usr/sbin/lircmd ] || { stop; exit 0; }

    if [ ! -f /etc/lirc/lircmd.conf ] || \
       grep -q "^#UNCONFIGURED" /etc/lirc/lircmd.conf
    then
        echo "Not configured; not starting."
        stop; exit 0
    fi
end script

exec /usr/sbin/lircmd
