# broken/memcached-autostart.conf
description "memcached-autostart - start memcached instances at boot"

start on filesystem

task

pre-start exec test -x /usr/bin/memcached || { stop; exit 0; }

script
    for i in /etc/memcached.d/*; do
        start -n memcached SERVER=${i#/etc/memcached.d/}
    done
end script
