# broken/memcached.conf
description "memcached - high-performance memory caching daemon"

start on filesystem
stop on runlevel [016]

expect daemon
respawn

env STARTER_SCRIPT=/usr/share/memcached/scripts/start-memcached
env ENABLE_MEMCACHED=yes

pre-start script
    [ -x /usr/bin/memcached -a -x $STARTER_SCRIPT ] || { stop; exit 0; }
    [ -s /etc/default/memcached ] && . /etc/default/memcached
    [ "$ENABLE_MEMCACHED" = "yes" ] || { stop; exit 0; }
end script

script
    NAME="memcached${SERVER:+_$SERVER}"
    $STARTER_SCRIPT /etc/${NAME}.conf /var/run/${NAME}.pid
end script
