# broken/restorecond.conf
description "restorecond - sets SELinux file context for watched directories"
author "Cameron Norman <camerontnorman@gmail.com>"

start on filesystem
stop on runlevel [016]

expect WHAT
respawn

pre-start script
	test -x /usr/sbin/restorecond || { stop; exit 0; }

	test -x /usr/sbin/selinuxenabled && /usr/sbin/selinuxenabled || {
		echo "SELinux disabled, not starting."
		stop; exit 0
	}

	test -f /etc/selinux/restorecond.conf || {
		echo "No restorecond config file found, not starting."
		stop; exit 0
	}
end script

exec /usr/sbin/restorecond
