# broken/usbmuxd.conf
description "usbmuxd - iOS device USB multiplexing daemon"

# I just need someone to test if this actually works.
# Also, the udev rule spawns its own daemon if systemd is
# not pid 1, so that needs to have an additional case for
# Upstart.

start on usb-device-added ID_VENDOR="05ac" ID_PRODUCT="/12[9a][0-9a-f]/*" 
stop on runlevel [016]

respawn

exec /usr/sbin/usbmuxd -f -u -U usbmux

pre-stop exec /usr/sbin/usbmuxd -X
