# broken/usbmuxd-try-stop.conf
description "usbmuxd-try-stop - only succeeds if the muxer has no devices"

# usbmuxd should probably just exit on idle, but hey

# See usbmuxd job
start on usb-device-added ID_VENDOR= ...

task

exec /usr/sbin/usbmuxd -x
