# broken/usbmuxd.conf
description "usbmuxd - iOS device USB multiplexing daemon"

# I just need someone to read the udev rule and convert it
# to upstart-udev-bridge speak. Look at configure-printer
# for an example. Also, use the upstart-monitor and plug in
# some usb device to get all that example stuff.
start on usb-device-added ID_VENDOR="05ac" ID_PRODUCT="/12[9a][0-9a-f]/*" 
stop on runlevel [016]

respawn

exec /usr/sbin/usbmuxd -f -u -U usbmux

pre-stop exec /usr/sbin/usbmuxd -X
