# broken/wicd.conf
description "wicd - network connection manager"

start on started dbus
stop on stopping dbus

expect daemon
respawn

pre-start script
    [ -x /usr/sbin/wicd ] || { stop; exit 0; }
    [ -s /etc/default/wicd ] && . /etc/default/wicd
    [ "$START_DAEMON" = "no" ] && { stop; exit 0; }
    [ -d /var/run/wicd ] || mkdir -p /var/run/wicd
end script

script
    [ -s /etc/default/wicd ] && . /etc/default/wicd
    exec /usr/sbin/wicd $DAEMON_ARGS
end script

post-stop exec rm -f /var/run/wicd/wicd.pid
