# display-manager/slim.conf
description "slim: Simple Login Manager"

start on (started dbus
          and multiuser-up
          and plymouth-ready)
stop on (stopped display-manager EXEC="slim")

emits login-session-start
emits desktop-session-start

env FORCE_START="no"

pre-start script
	[ -n "$UPSTART_EVENTS" ] || FORCE_START="yes"
	exec start display-manager EXEC="slim" FORCE="$FORCE_START"
end script

post-stop exec initctl stop display-manager EXEC="slim"
