# database/firebird2.5-classic.conf
description "firebrid2.5-classic - firebird database inetd prep"

start on filesystem

env FB_VER="2.5"
env FB_FLAVOUR="classic"

pre-start script
    test -x /usr/sbin/fb_inet_server || { stop; exit 0; }
    FUNCTIONS=/usr/share/firebird${FB_VER}-common/functions.sh
    test -e $FUNCTIONS || { stop; exit 0; }
    . $FUNCTIONS
    create_var_run_firebird
end script
