# database/firebird3.0-rundir.conf
start on starting firebird3.0
task
exec install -o firebird -g firebird -d /run/firebird3.0
