# database/nutcracker.conf
description "memcached and redis proxy"

start on (runlevel [2345] and net-device-up IFACE!=lo)
stop on runlevel [!2345]

respawn

pre-start script
  [ ! -e /etc/nutcracker/nutcracker.yml ] && { stop; exit 0; }

  mkdir -p -m0755 /var/run/nutcracker
  chown nutcracker:nutcracker /var/run/nutcracker

  . /etc/default/nutcracker
  /usr/sbin/nutcracker $DAEMON_OPTS --test-conf || { stop; exit 1; }
end script

script
  . /etc/default/nutcracker
  exec start-stop-daemon --quiet --start --chuid nutcracker \
       --exec /usr/sbin/nutcracker -- $DAEMON_OPTS
end script
