# database/redis-server.conf
description "redis-server - Persistent key-value database"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork

env CONF=/etc/redis/redis.conf
env RUNDIR=/var/run/redis
env USER=redis

umask 007
# Set soft and hard to appropiate values for high concurrency environments
#limit nofile soft hard

pre-start script
    [ -x /usr/bin/redis-server ] || { stop; exit 0; }    
    [ -d $RUNDIR ] || install -o $USER -g $USER -m 755 -d $RUNDIR
    touch $RUNDIR/redis-server.pid
end script

exec start-stop-daemon -S -q -c $USER -x /usr/bin/redis-server -- $CONF

post-stop exec rm -f $RUNDIR/redis-server.pid
