# database/tarantool.conf
description "tarantool - NoSQL Database"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

pre-start script
    SUCCESS=true
    INSTANCES=$(find /etc/tarantool/instances.enabled -xtype f -name '*.cfg')
    for inst in $INSTANCES; do
        tarantool_instance $inst start || SUCCESS=false
    done
    test "$SUCCESS" = "true" || exit 1
end script

post-stop script
    SUCCESS=true
    INSTANCES=$(find /etc/tarantool/instances.enabled -xtype f -name '*.cfg')
    for inst in $INSTANCES; do
        tarantool_instance $inst stop || SUCCESS=false
    done
    test "$SUCCESS" = "true" || exit 1
end script