# database/yrmcds.conf
description "yrmcds: memcached compatible key value store"

start on net-device-up
stop on runlevel [!2345]

limit nofile 100000 100000
umask 077
chdir /tmp
setuid yrmcds
setgid yrmcds

respawn
exec /usr/sbin/yrmcdsd
