# database/yrmcds-prep.conf
description "yrmcds-prep: memcached compatible key value store (prep)"

start on starting yrmcds
stop on stopped yrmcds

pre-start script
	install -m 0700 -o yrmcds -g yrmcds -d /var/tmp/yrmcds
	rm -f /var/tmp/yrmcds/*

	# make sure keepalived init script was disabled.
	# ex) sudo update-rc.d keepalived disable
	#/etc/init.d/keepalived start
	#sleep 5
end script

#post-stop script
#    /etc/init.d/keepalived stop
#
#    # sleep a while to hand-off VIP
#    sleep 5
#end script
