# rc/rc.conf
description "rc: sysvinit runlevel compatibility"

start on runlevel [0123456]
stop on runlevel [!$RUNLEVEL]

export RUNLEVEL
export PREVLEVEL

task
console output
exec /etc/init.d/rc $RUNLEVEL
