# rc/rc.local.conf
description "rc.local: run /etc/rc.local after multiuser runlevel"

start on stopped rc RUNLEVEL=[2345]

pre-start script
    [ -x /etc/rc.local ] || { stop; exit 0; }
    exec rc.local
end script
