# rc/rc.local.conf
description "rc.local: run /etc/rc.local after multiuser runlevel"

start on stopped rc RUNLEVEL=[2345]

pre-start script
    test -x /etc/rc.local || exit 0
    exec /etc/rc.local
end script
