# fs/filesystem.conf
description "filesystem: synchronization job for fs checkpoint"

start on fs-up and filesystem
