# fs/umountall.conf
description "unmount all filesystems"

start on fs-down

emits unmounting-filesystem
emits unmounted-remote-filesystems

exec /sbin/umountall
