# INSTALL.minimal/container-ttys.conf
description "container-ttys: start getty jobs in a container"

start on startup

env container_ttys

pre-start script
	for t in $container_ttys; do
		start -q -n getty PORT=${t#/dev/} || true
	done
end script
