# INSTALL.minimal/control-alt-delete.conf
description	"control-alt-delete: emergency keypress handling"
author		"Scott James Remnant <scott@netsplit.com>"

# special event emitted by Upstart to indicate when
# the Control-Alt-Delete key combination is pressed
start on control-alt-delete

task
exec shutdown -r now "Control-Alt-Delete pressed"
