# INSTALL.minimal/firewall.conf
description "firewall: firewall checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on starting net

emits firewall-up
emits firewall-down

pre-start exec initctl emit firewall-up
post-stop exec initctl emit firewall-down
