# INSTALL.minimal/graphical.conf
description "graphical: graphical checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on started multiuser
stop on stopping multiuser

emits graphical-up
emits graphical-down

pre-start exec initctl emit graphical-up
post-stop exec initctl emit graphical-down
