# INSTALL.minimal/hostname.conf
description "hostname: set hostname from /etc/hostname"

start on sysinit-up

pre-start exec hostname -b -F /etc/hostname
