# INSTALL.minimal/sysinit.conf
description "sysinit: system initialization checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on startup
stop on runlevel [06]

emits sysinit-up
emits sysinit-down

pre-start exec initctl emit sysinit-up
post-stop exec initctl emit sysinit-down
