# INSTALL.minimal/syslog.conf
description "syslog: system logging checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on basic-up
stop on basic-down

emits syslog-up
emits syslog-down

pre-start exec initctl emit syslog-up
post-stop exec initctl emit syslog-down
