# fs/mountall/mountall-net.conf
description "mountall-net: mount network filesystems"
start on net-device-up
task
console none
exec reload mountall
