# fs/mountall/mountall-shell.conf
description "mountall-shell: recovery shell for filesystem failure"

start on (stopped mountall EXIT_STATUS=[!4]
          or stopped mountall EXIT_SIGNAL=?*)
stop on runlevel [06]

task
console owner

script
    action="reboot the system."
    case "$EXIT_STATUS" in
    ""|1)
    	err="General error mounting filesystems."
	;;
    2)
    	err="Filesystem check or mount failed."
	action="continue boot,\nre-trying filesystems and ignoring further errors."
	;;
    3)
        err="Root filesystem check failed."
	;;
    esac

    echo "${err}"
    echo "A maintenance shell will now be started."
    echo "CONTROL-D will terminate this shell and ${action}"

    /sbin/sulogin
end script

post-stop script
    [ -z "$UPSTART_STOP_EVENTS" ] || exit 0
    [ "$EXIT_STATUS" != "2" ] || { start -n mountall; exit 0; }
    umount -a || :
    exec reboot -f
end script
