# broken/lirc/lirc-irexec.conf
description "lirc-irexec - "

start on starting lirc and started lirc-lircd
stop on stopping lirc

expect fork
respawn

env CONF=/etc/lirc/lircrc

pre-start script
    [ -x /usr/bin/irexec ] || { stop; exit 0; }

    if [ ! -f $CONF ] || \
       grep -q "^#UNCONFIGURED" $CONF
    then
        echo "Not configured; not starting."
        stop; exit 0
    fi
end script

exec /usr/bin/irexec -d $CONF
