# broken/lirc/lirc.conf
description "lirc - "

start on filesystem
stop on runlevel [016]

# This is a meta job.
# All the LIRC daemons hook into this job's start and stop.
