# net/netscript/netscript.conf
# netscript - configure network devices controlled by netscript
#
#  
#
# Based on work by 
# Scott James Remnant <scott@ubuntu.com>  Tue, 15 Sep 2009 03:30:29 +0100
# copyright 2009 Canonical Limited. License GPL v2
# Adaptation for netscript-2.4
# copyright 2011 Matthew Alexander Grant <matthewgrant5@gmail.com>

description	"configure virtual network devices"

start on (local-filesystems
	  and (stopped udevtrigger or container))

task

exec netscript start
