# net/netscript/netscript-interface.conf
# network-interface - configure network device
#
# This service causes network devices to be brought up or down as a result
# of hardware being added or removed, including that which isn't ordinarily
# removable.
#
# Based on work by 
# Scott James Remnant <scott@ubuntu.com>  Tue, 15 Sep 2009 03:30:29 +0100
# copyright 2009 Canonical Limited. License GPL v2
# changes for netscript-2.4
# copyright 2011 Matthew Alexander Grant License GPL V2
# Matthew Grant <matthewgrant5@gmail.com>

description	"configure network device"

start on stopped rc RUNLEVEL=[2345] and net-device-added
stop on runlevel [2345] and net-device-removed INTERFACE=$INTERFACE

instance $INTERFACE

pre-start script
	if [ "$INTERFACE" != "lo" ]; then
		exec netscript ifup $INTERFACE
	fi
end script

post-stop script
	if [ "$INTERFACE" != "lo" ]; then
		exec netscript ifdown $INTERFACE
	fi
end script
