# net/netscript/netscript-shutdown.conf
# netscript shutdown - shutdown system networking controlled by netscript
#
# This task stops all the network daemons and interfaces controlled by
# netscript

description	"shutdown system networking controlled by netscript"

start on runlevel [06]

task

exec netscript stop
