# httpd/htcacheclean.conf
description "htcacheclean - HTTP cache cleaning daemon for Apache"

start on started apache2
stop on stopped apache2 AP_INSTANCE=$AP_INSTANCE

instance $AP_INSTANCE

expect fork

# Overridden in /etc/default/apache2{,-xxx}
env HTCACHECLEAN_RUN=auto
env HTCACHECLEAN_MODE=daemon
env HTCACHECLEAN_SIZE=300M
env HTCACHECLEAN_DAEMON_INTERVAL=120
env HTCACHECLEAN_PATH=/var/cache/apache2/mod_cache_disk
env HTCACHECLEAN_OPTIONS=

pre-start script
    [ -x /usr/bin/htcacheclean ] || { stop; exit 0; }

    [ -n $AP_INSTANCE ] && suffix="-$AP_INSTANCE"

    if [ -f /etc/default/apache2$suffix ]; then
        . /etc/default/apache2$suffix
    elif [ -f /etc/default/apache2 ]; then
        . /etc/default/apache2
    fi

    [ "$HTCACHECLEAN_MODE" = "daemon" ] || { stop; exit 0; }
    
    # Must be yes, or auto + module is enabled, for this job to proceed
    if [ "$HTCACHECLEAN_RUN" = "yes" ] ; then exit 0; fi
    if [ "$HTCACHECLEAN_RUN" = "auto" ]; then
        APACHE_CONFDIR=/etc/apache2$suffix
        if [ -r $APACHE_CONFDIR/envvars ] ; then
            . $APACHE_CONFDIR/envvars
        fi 
        if [ -e ${APACHE_MODS_ENABLED:-$APACHE_CONFDIR/mods-enabled}/cache_disk.load ]; then
            exit 0
        fi
    fi

    stop
    exit 0
end script

script
    [ -n $AP_INSTANCE ] && suffix="-$AP_INSTANCE"

    # We need the suffix for this one
    HTCACHECLEAN_PATH=/var/cache/apache2${suffix}/mod_cache_disk

    if [ -f /etc/default/apache2$suffix ]; then
        . /etc/default/apache2$suffix
    elif [ -f /etc/default/apache2 ]; then
        . /etc/default/apache2
    fi

    exec htcacheclean $HTCACHECLEAN_OPTIONS -d$HTCACHECLEAN_DAEMON_INTERVAL \
                      -i -p$HTCACHECLEAN_PATH -l$HTCACHECLEAN_SIZE
end script
