# proxy/squid.conf
description "squid - proxy caching server"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem

# TODO: readiness
respawn

limit nofile 1024 1024
umask 027

exec /usr/sbin/squid -N -sDYC
