# grandfathered/avahi-daemon.conf
description "Avahi - mDNS/DNS-SD Zeroconf daemon"

apparmor load /etc/apparmor.d/usr.sbin.avahi-daemon

start on started dbus
stop on stopping dbus or runlevel [016] or unmounting-filesystem

expect daemon
respawn

exec /usr/sbin/avahi-daemon -D
