# grandfathered/cron.conf
description "cron - peridoic command scheduler"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect daemon
respawn

script
    [ -r /etc/default/cron ] && . /etc/default/cron
    
    if [ "$READ_ENV" = "yes" ] ; then
        for ENV_FILE in /etc/environment /etc/default/locale; do
            [ -r "$ENV_FILE" ] || continue
            [ -s "$ENV_FILE" ] || continue

            for var in LANG LANGUAGE LC_ALL LC_CTYPE; do
                value=`egrep "^${var}=" "$ENV_FILE" | tail -n1 | cut -d= -f2`
                [ -n "$value" ] && eval export $var=$value
            done
         done

        # Get the timezone set.
        if [ -z "$TZ" -a -e /etc/timezone ]; then
            TZ=`cat /etc/timezone` 
        fi
    fi

    exec cron -L 8
end script
