# grandfathered/cups.conf
description "CUPS - Common Unix Printing Service"

apparmor load /etc/apparmor.d/usr.sbin.cupsd

start on filesystem
stop on runlevel [016] or unmounting-filesystem

expect fork
respawn
respawn limit 3 12

pre-start script
    [ -x /usr/sbin/cupsd ] || { stop; exit 0; }
    mkdir -p /var/run/cups/certs
end script

exec /usr/sbin/cupsd
