# grandfathered/network-manager.conf
description "NetworkManager - automatically connect to wired and wireless networks"

start on started dbus and static-network-up
stop on stopping dbus or runlevel [016]

expect fork

exec NetworkManager
