# grandfathered/ufw.conf
description "ufw - Uncomplicated FireWall"

start on (starting network-interface
          or starting network-manager
          or starting networking)
stop on runlevel [1]

pre-start exec /lib/ufw/ufw-init start quiet
post-stop exec /lib/ufw/ufw-init stop
