# time/chrony.conf
description "chrony: an NTP client/server"
#TODO: periodic (hourly) dns srv lookup using /usr/lib/chrony/chrony-helper update-dnssrv-servers %I

start on started sysinit
stop on stopping sysinit

env CMD=/usr/sbin/chronyd
env SYSD_CMD=/usr/lib/systemd/scripts/chronyd-starter.sh

pre-start script
	test -x "${CMD}" || { stop; exit 0; }
end script

script
	test ! -x "$SYSD_CMD" || CMD="$SYSD_CMD"
	exec "$CMD" -n
end script

post-start exec /usr/lib/chrony/chrony-helper update-daemon || true
