# time/ntpsec.conf
description "NTPsec: secure, hardened, and improved Network Time Protocol daemon"

# TODO: /usr/lib/ntp/rotate-stats daily

start on started nss-lookup
stop on stopping nss-lookup

apparmor load /etc/apparmor.d/usr.sbin.ntpd

env RUNASUSER=ntp
env NTPD_OPTS="-n -g"

# Specifying -g on the command line allows ntpd to make large adjustments to
# the clock on boot.  However, if `respawn` is used, a malicious (or broken)
# server could send the incorrect time, trip the panic threshold, and when
# ntpd restarts, serve it the incorrect time (which would be accepted).

script
    [ ! -e /run/ntp.conf.dhcp ] || CONF="-c /run/ntp.conf.dhcp"

    UGID=$(getent passwd $RUNASUSER | cut -f 3,4 -d:) || true
    if test "$(uname -s)" = "Linux" && test -n "$UGID"; then
            NTPD_OPTS="$NTPD_OPTS -u $UGID"
    fi

    (
    if flock -w 180 9; then
        exec /usr/sbin/ntpd ${NTPD_OPTS} ${CONF}
    else
        echo "Timeout waiting for $LOCKFILE"
        exit 1
    fi
    ) 9>/run/lock/ntpdate
end script

post-start exec /usr/sbin/ntpwait -s 1 -n 30000
