# time/openntpd.conf
description "openntpd: OpenBSD Network Time Protocol daemon"

start on started sysinit
stop on stopping sysinit

apparmor load /etc/apparmor.d/usr.sbin.ntpd

pre-start script
	test ! -s /etc/default/openntpd || . /etc/default/openntpd
	test -d /run/openntpd || mkdir -p /run/openntpd
	exec /usr/sbin/ntpd -n $DAEMON_OPTS
end script

respawn
script
	test ! -s /etc/default/openntpd || . /etc/default/openntpd
	exec /usr/sbin/ntpd -d $DAEMON_OPTS
end script
