# selinux/selinux-mounted-varlock.conf
# selinux-mounted-varlock - Run restorecon on /var/lock directory

description	"Run restorecon on newly mounted /var/lock directory"

start on mounted MOUNTPOINT=/var/lock
env MOUNTPOINT=/var/lock

task

script
    if [ -x /sbin/restorecon ]; then
        /sbin/restorecon -R "${MOUNTPOINT}"
    fi
end script
