# selinux/selinux-mounted-varrun.conf
# selinux-mounted-varrun - Run restorecon on /var/run directory

description	"Run restorecon on newly mounted /var/run directory"

start on mounted MOUNTPOINT=/var/run
env MOUNTPOINT=/var/run

task

script
    if [ -x /sbin/restorecon ]; then
        /sbin/restorecon -R "${MOUNTPOINT}"
    fi
end script
