# selinux/selinux-mounted-dev.conf
# selinux-mounted-dev - Run restorecon on /dev directory

description	"Run restorecon on newly mounted /dev directory"

start on mounted MOUNTPOINT=/dev
env MOUNTPOINT=/dev

task

script
    if [ -x /sbin/restorecon ]; then
        /sbin/restorecon -R "${MOUNTPOINT}"
    fi
end script
