# selinux/selinux-mounted-tmp.conf
# selinux-mounted-tmp - Run restorecon on /tmp directory

description	"Run restorecon on newly mounted /tmp directory"

start on mounted MOUNTPOINT=/tmp
env MOUNTPOINT=/tmp

task

script
    if [ -x /sbin/restorecon ]; then
        /sbin/restorecon "${MOUNTPOINT}"
    fi
end script
