# base/upstart/halt.conf
description "halt"

start on runlevel [0] and stopped sysinit

# By default, we use the environment of the runlevel event
# to decide whether to halt or poweroff.
# Setting this variable to HALT or POWEROFF means
# the environment is ignored and the respective action is always taken.
env INIT_HALT=""

# uncomment to enable wake-on-lan
#env NETDOWN=no

script
	if [ "$INIT_HALT" = "" ]; then
		case "$HALT" in
		  [Hh]*) INIT_HALT=HALT ;;
		esac
	fi

	test "$INIT_HALT" = "HALT" || {
		poweroff="-p"
		[ -x /etc/init.d/ups-monitor ] && /etc/init.d/ups-monitor poweroff || :
		[ -x /sbin/initctl ] && /sbin/initctl emit --quiet ups-poweroff 2>/dev/null || :
	}

	# Don't shut down drives if we're using RAID
	grep -qs '^md.*active' /proc/mdstat || hddown="-h"

	# Don't shut down net interfaces if we're using wake-on-lan
	test "$NETDOWN" = "no" || netdown="-i"

	exec halt -d -f $netdown $poweroff $hddown
end script
