# base/upstart/power-status-changed.conf
description "power-status-changed: trigger immediate shutdown on SIGPWR"

# special event emitted by Upstart when it receives SIGPWR
# this occurs when certain UPS tools or containers request a shutdown
start on power-status-changed

task
exec shutdown -h now "SIGPWR received"
