# base/upstart/rc-sysinit.conf
description	"rc-sysinit: System V initialisation compatibility"
author		"Scott James Remnant <scott@netsplit.com>"

start on started multiuser or failsafe-boot
stop on runlevel
emits runlevel

# There can be no previous runlevel here, but there might be old
# information in /var/run/utmp that we pick up, and we don't want
# that, so we override it.
env RUNLEVEL=
env PREVLEVEL=
env DEFAULT_RUNLEVEL=2
env INIT_VERBOSE

console output
task
exec telinit "${DEFAULT_RUNLEVEL}"
