# base/upstart/reboot.conf
description "reboot"
start on runlevel [6] and stopped sysinit
exec reboot -d -f -i
