# base/upstart/reboot.conf
description "reboot"
start on runlevel [6]
exec reboot -d -f -i
