# base/upstart/sysinit.conf
description "sysinit: checkpoint for initial setup / multi user services"

start on filesystem and static-network-up
stop on unmounting-filesystem or deconfiguring-networking
