# base/upstart/upstart-file-bridge.conf
description "upstart-file-bridge"

# the bridge does not currently handle mounts
# that overlay already-watched directories,
# so we must start when the fs is already up

start on basic-up
stop on basic-down
emits file

expect daemon
respawn
exec upstart-file-bridge --daemon
