# base/upstart/upstart-socket-bridge.conf
description "upstart-socket-bridge"

start on basic-up
stop on basic-down
emits socket

expect daemon
respawn
exec upstart-socket-bridge --daemon
