# checkpoints/basic.conf
description "basic: pre multiuser checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on started sysinit
stop on stopping sysinit

emits basic-up
emits basic-down

pre-start exec initctl emit basic-up
post-stop exec initctl emit basic-down
