# checkpoints/failsafe.conf
description "failsafe: fallback init sequence"

# wait for all the events between
# local-fs/coldplug and the basic checkpoint,
# (most importantly static-network-up)
start on local-filesystems and coldplug-trigger
stop on started basic
emits failsafe-boot

env TIMEOUT=120
script
	sleep "${TIMEOUT}"
	initctl emit -n failsafe-boot
        logger -t 'failsafe' -p daemon.warning "Failsafe of ${TIMEOUT} seconds reached."
end script
