# checkpoints/multiuser.conf
description "multiuser: multiuser checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on started basic
stop on stopping basic

emits multiuser-up
emits multiuser-down

pre-start exec initctl emit multiuser-up
post-stop exec initctl emit multiuser-down
