# checkpoints/net.conf
description "net: static networking services checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on sysinit-up and local-filesystems and coldplug-trigger
stop on unmounted-remote-filesystems

emits net-up
emits net-down

pre-start exec initctl emit net-up
post-stop exec initctl emit net-down
