# checkpoints/nss-lookup.conf
description "nss-lookup: nss-lookup checkpoint"
author "Cameron Nemo <cnemo@protonmail.com>"

start on basic-up and started syslog
stop on stopping syslog

emits nss-lookup-up
emits nss-lookup-down

pre-start exec initctl emit nss-lookup-up
post-stop exec initctl emit nss-lookup-down
