# net/isc-dhcp/isc-dhcp-relay.conf
description "ISC DHCP IPv4 relay"
author "Stéphane Graber <stgraber@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
    if [ ! -f /etc/default/isc-dhcp-relay ]; then
        echo "/etc/default/isc-dhcp-relay does not exist! - Aborting..."
        echo "Run 'dpkg-reconfigure isc-dhcp-relay' to fix the problem."
        stop
        exit 0
    fi
end script

respawn
script
    . /etc/default/isc-dhcp-relay

    # Build command line for interfaces (will be passed to dhrelay below.)
    IFCMD=""
    if test "$INTERFACES" != ""; then
        for I in $INTERFACES; do
            IFCMD=${IFCMD}"-i "${I}" "
        done
    fi

    exec /usr/sbin/dhcrelay -d -q -4 $OPTIONS $IFCMD $SERVERS
end script
