# net/isc-dhcp/isc-dhcp-relay6.conf
description "ISC DHCP IPv6 relay"
author "Stéphane Graber <stgraber@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
    if [ ! -f /etc/default/isc-dhcp-relay6 ]; then
        echo "/etc/default/isc-dhcp-relay6 does not exist! - Aborting..."
        stop
        exit 0
    fi

    . /etc/default/isc-dhcp-relay6
    if [ -z "$LOWER_INTERFACES" ] || [ -z "$UPPER_INTERFACES" ]; then
        echo "No interface defined in /etc/default/isc-dhcp-relay6! - Aborting..."
        stop
        exit 0
    fi
end script

respawn
script
    . /etc/default/isc-dhcp-relay6

    # Build command line for interfaces (will be passed to dhrelay below.)
    IFCMD=""
    if test "$LOWER_INTERFACES" != ""; then
        for I in $LOWER_INTERFACES; do
            IFCMD=${IFCMD}"-l "${I}" "
        done
    fi
    if test "$UPPER_INTERFACES" != ""; then
        for I in $UPPER_INTERFACES; do
            IFCMD=${IFCMD}"-u "${I}" "
        done
    fi

    exec /usr/sbin/dhcrelay -d -q -6 $OPTIONS $IFCMD
end script
