# syslog/busybox-klogd.conf
description "busybox-klogd - Busybox kernel logger"

start on syslog-up
stop on syslog-down

expect fork
respawn
reload signal 0

pre-start script
    test -x /sbin/klogd || { stop; exit 0; }
end script

exec /sbin/klogd
