# syslog/nxlog.conf
description "nxlog"

start on syslog-up
stop on syslog-down

pre-start exec install -o nxlog -g nxlog -m 1770 -d /var/run/nxlog

expect #TODO
exec /usr/bin/nxlog
