# syslog/rsyslog.conf
description "rsyslog: enhanced, multi threaded syslogd"

start on starting sysinit
stop on stopped sysinit

apparmor load /etc/apparmor.d/usr.sbin.rsyslogd

pre-start script
	test -x /usr/sbin/rsyslogd || { stop; exit 0; }
end script

expect fork
respawn
exec /usr/sbin/rsyslogd
