# syslog/ulogd2.conf
# ulogd2 - Netfilter Userspace Logging Daemon
#
# ulogd is an advanced netfilter logging daemon.

description	"Netfilter Userspace Logging Daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

exec /usr/sbin/ulogd --daemon --uid ulog
