# base/anacron.conf
description "anac(h)ronistic cron"

start on started sysinit
stop on stopping sysinit

expect fork
normal exit 0

exec anacron -s
