# base/cron.conf
description "cron: periodic command scheduler"

start on filesystem
stop on unmounting-filesystem

respawn
exec cron -f
