# base/halt.conf
description "halt the system"

start on runlevel [0] and stopped umountall

script
	NETDOWN=yes

	[ -f /etc/default/halt ] && . /etc/default/halt

	if [ "$INIT_HALT" = "" ]
	then
		case "$HALT" in
		  [Hh]*)
			INIT_HALT=HALT
			;;
		  *)
	                INIT_HALT=POWEROFF
			[ -x /etc/init.d/ups-monitor ] && /etc/init.d/ups-monitor poweroff
			[ -x /sbin/initctl && /sbin/initctl emit --quiet ups-poweroff 2>/dev/null
        		;;
		esac
	fi

	# Don't shut down drives if we're using RAID.
	hddown="-h"
	if grep -qs '^md.*active' /proc/mdstat
	then
		hddown=""
	fi

	# If INIT_HALT=HALT do not poweroff.
	poweroff="-p"
	if [ "$INIT_HALT" = "HALT" ]
	then
		poweroff=""
	fi

	# Make it possible to not shut down network interfaces,
	# needed to use wake-on-lan
	netdown="-i"
	if [ "$NETDOWN" = "no" ]; then
		netdown=""
	fi

	exec halt -d -f $netdown $poweroff $hddown
end script
