# base/hostname.sh.conf
description "Set the hostname"

start on stopped hostname
