# base/motd.conf
description "motd - dynamic creation of /run/motd.dynamic"

start on started hostname and virtual-filesystems

pre-start exec uname -snrvm > /run/motd.dynamic
