# base/nscd.conf
description "nscd: name service cache daemon"

start on filesystem
stop on unmounting-filesystem

umask 022

pre-start script
    [ -x /usr/sbin/nscd ] || { stop; exit 0; }
    [ -s /etc/nscd.conf ] || { stop; exit 0; }
    [ -d /var/run/nscd ] || mkdir -p /var/run/nscd
end script

reload signal 0
expect fork
respawn
exec /usr/sbin/nscd

pre-stop exec /usr/sbin/nscd --shutdown

post-stop exec rm -f /var/run/nscd/nscd.pid
