# base/procps.conf
description "procps - Configure kernel parameters at boot"

start on starting network-pre

pre-start exec /sbin/sysctl --system -q
