# base/procps.conf
description "procps - Configure kernel parameters at boot"

start on (virtual-filesystems or static-network-up)

pre-start exec /sbin/sysctl --system -q
