# base/reboot.conf
description "reboot the system"

start on runlevel [6]

exec reboot -d -f -i
