# base/reboot.conf
description "reboot the system"

start on runlevel [6] and stopped umountall

exec reboot -d -f -i
