# base/rsyslog.conf
description "rsyslog: enhanced, multi threaded syslogd"

start on filesystem
stop on unmounting-filesystem

apparmor load /etc/apparmor.d/usr.sbin.rsyslogd

expect fork
respawn
exec /usr/sbin/rsyslogd
