# base/umountall.conf
description "unmount all filesystems"

start on starting reboot or starting halt

emits unmounting-filesystem
emits unmounted-remote-filesystems

exec /sbin/umountall
