# base/umountall.conf
description "unmount all filesystems"

start on stopped rc RUNLEVEL=[06]

emits unmounting-filesystem
emits unmounted-remote-filesystems

exec /sbin/umountall
