# mail/dovecot.conf
description "dovecot - secure and highly configurable mail server"

start on runlevel [2345]
stop on runlevel [!2345]

env CONF=/etc/dovecot/dovecot.conf

# Uncomment the following line to allow Dovecot to produce core dumps
#limit nofile unlimited unlimited

pre-start script
    [ -x /usr/sbin/dovecot ] || { stop; exit 0; }
    [ -s $CONF ] || { stop; exit 0; }
    # TODO: after ntp
end script

respawn
exec /usr/sbin/dovecot -F -c $CONF

apparmor load /etc/apparmor.d/usr.sbin.dovecot
apparmor load /etc/apparmor.d/usr.lib.dovecot.deliver
apparmor load /etc/apparmor.d/usr.lib.dovecot.dovecot-auth
apparmor load /etc/apparmor.d/usr.lib.dovecot.imap
apparmor load /etc/apparmor.d/usr.lib.dovecot.imap-login
apparmor load /etc/apparmor.d/usr.lib.dovecot.managesieve-login
apparmor load /etc/apparmor.d/usr.lib.dovecot.pop3
apparmor load /etc/apparmor.d/usr.lib.dovecot.pop3-login
