# mail/exim4.conf
description "exim4 - Mail Transfer Agent"

start on filesystem and static-network-up
stop on runlevel [016] or unmounting-filesystem

env RUNDIR=/var/run/exim4

pre-start script
    [ -x /usr/sbin/exim4 ] || { stop; exit 0; }
 
    UPEX4CONF=$(command -v update-exim4.conf)
    [ -z $UPEX4CONF ] || $UPEX4CONF

    # Stop this job if exim is run from inet
    if [ -f /etc/inetd.conf ] && grep -E -q '^[[:space:]]*((\*|[[:alnum:].-]+):)?smtp[[:space:]]' /etc/inetd.conf
    then
        stop; exit 0
    fi

    [ ! -s /etc/default/exim4 ] || . /etc/default/exim4 || :
    [ "$QUEUERUNNER" != "nodaemon" ] || { stop; exit 0; }

    # Create runtime directory if not present
    install -d -m 750 -o Debian-exim -g Debian-exim $RUNDIR

    # emit blocking event that causes this job to wait for
    # the listener and queue to start
    initctl emit start-exim
end script

post-stop script
    # emit blocking event for the stop
    initctl emit stop-exim
    rm -f $RUNDIR/exim.pid
end script
