# mail/postfix.conf
description "postfix: mail transport agent"

start on started nss-lookup
stop on stopping nss-lookup

pre-start script
	test -x /usr/sbin/postfix && test -f /etc/postfix/main.cf || { stop; exit 0; }
	exec >/dev/null 2>&1
	for i in $(postmulti -l -a | awk '($3=="y") { print $1}'); do
		start postfix-instance INSTANCE="$1" || :
	done
end script
