# radosgw/radosgw-agent.conf
# vim: set ft=upstart ts=2 et:
description "RADOS Gateway Sync Agent"
author "James Page <james.page@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
  ENABLED="no"
  if [ -f /etc/default/radosgw-agent ]; then
    . /etc/default/radosgw-agent
  fi
  if [ "x$ENABLED" = "xyes" ]; then
    exec /usr/bin/radosgw-agent -c /etc/ceph/radosgw-sync.conf \
      --log-file /var/log/radosgw/radosgw-agent.log
  fi
end script
