# fso/phonefsod.conf
description "phonefsod - SHR Phone freesmartphone.org daemon"

start on started dbus
stop on stopping dbus

expect fork
respawn

# Daemon can not reload
reload signal 0

pre-start script
    test -x /usr/sbin/phonefsod || { stop; exit 0; }
end script

exec /usr/sbin/phonefsod
