# net/avahi/avahi-dnsconf.conf
description "avahi-dnsconfd - Avahi Unicast DNS Configuration Daemon"

start on started avahi-daemon
stop on stopping avahi-daemon or runlevel [016] or unmounting-filesystem

expect daemon
respawn

pre-start script
    if [ -s /etc/localtime ]; then
        [ -d /etc/avahi/etc ] || mkdir -p /etc/avahi/etc
        cp -fp /etc/localtime /etc/avahi/etc
    fi
end script

exec avahi-dnsconfd -D
