# net/avahi/avahi-cups-reload.conf
description "Reload cups, upon starting avahi-daemon to make sure remote queues are populated"
author "Dmitrijs Ledkovs <dmitrijs.ledkovs@canonical.com>"

start on started avahi-daemon
task
exec reload cups
