# nss-lookup/bind9-pkcs11.conf
description "bind9-pkcs11: BIND DNS w/ native PKCS#11"

start on starting nss-lookup
stop on stopped nss-lookup

apparmor load /etc/apparmor.d/usr.sbin.named

env KRB5_KTNAME=/etc/bind/named.keytab

pre-start exec install -u root -g bind -m 775 -d /var/run/named
exec /usr/sbin/named-pkcs11 -f -u bind
pre-stop exec rndc stop
