# nss-lookup/bind9.conf
description "bind9: BIND Domain Name Server"

start on starting nss-lookup
stop on stopped nss-lookup

apparmor load /etc/apparmor.d/usr.sbin.named

pre-start exec install -u root -g bind -m 775 /var/run/named

exec /usr/sbin/named -f -u bind

post-start script
    [ -x /sbin/resolvconf ] || exit 0
    echo nameserver 127.0.0.1 | /sbin/resolvconf -a lo.named
end script

pre-stop exec rndc stop

post-stop script
    [ -x /sbin/resolvconf ] || exit 0
    /sbin/resolvconf -d lo.named
end script
