# nss-lookup/knot.conf
# knot - Knot DNS Server

description "Knot DNS Server"
author "Ondřej Surý <ondrej@debian.org>"

start on nss-lookup-up
stop on nss-lookup-up

env CONFFILE=/etc/knot/knot.conf

respawn
pre-start exec /usr/lib/knot/prepare-environment $CONFFILE
exec /usr/sbin/knotd -c $CONFFILE
