# nss-lookup/nsd.conf
# nsd - Name Server Daemon

description "Name Server Daemon"
author "Ondřej Surý <ondrej@debian.org>"

start on nss-lookup-up
stop on nss-lookup-up

env CONFFILE=/etc/nsd/nsd.conf

respawn
pre-start script
    # prepare environment
    mkdir -p "$(dirname "$(/usr/sbin/nsd-checkconf -o pidfile $CONFFILE)")"
    chown "$(/usr/sbin/nsd-checkconf -o username $CONFFILE)" "$(dirname "$(/usr/sbin/nsd-checkconf -o pidfile $CONFFILE)")"
    mkdir -p "$(dirname "$(/usr/sbin/nsd-checkconf -o database $CONFFILE)")"
    chown "$(/usr/sbin/nsd-checkconf -o username $CONFFILE)" "$(dirname "$(/usr/sbin/nsd-checkconf -o database $CONFFILE)")"
end script

exec /usr/sbin/nsd -d -c $CONFFILE
