# nss-lookup/nss-lookup.conf
description "nss-lookup: DNS server checkpoint"

start on started sysinit
stop on stopping sysinit
