# virt/lxc/lxc-net.conf
description "lxc network"
author "Serge Hallyn <serge.hallyn@canonical.com>"

start on starting lxc
stop on stopped lxc

pre-start exec /usr/lib/arm-linux-gnueabihf/lxc/lxc-net start
post-stop exec /usr/lib/arm-linux-gnueabihf/lxc/lxc-net stop
