# virt/lxc/lxc-net.conf
description "lxc network"
author "Serge Hallyn <serge.hallyn@canonical.com>"

start on starting lxc-autostart
stop on stopped lxc-autostart

pre-start exec /usr/lib/lxc/lxc-net start
post-stop exec /usr/lib/lxc/lxc-net stop
