# net/dhclient-container.conf
description "dhclient: acquire DHCP leases"
start on container CONTAINER=lxc
manual
exec dhclient -d -q
pre-stop exec dhclient -r
