# net/dhclient.conf
description "dhclient: acquire DHCP leases"

# probably need some device matching here
start on net-device-added or container CONTAINER=lxc
stop on net-device-removed INTERFACE=$INTERFACE

# this will clash with any other network configuration system
manual

env INTERFACE=""
instance $INTERFACE

# dhclient.leases database
env PATH_DHCLIENT_DB=/var/lib/dhcp/dhclient.${INTERFACE:+$INTERFACE.}leases
# dhclient.pid file
env PATH_DHCLIENT_PID=/run/dhclient.${INTERFACE:+$INTERFACE.}pid

exec /sbin/dhclient -d -q $INTERFACE
pre-stop exec /sbin/dhclient -r
