# net/lxcguestnet.conf
description "lxcguestnet: minimal networking in Linux Containers"

start on net-up and container CONTAINER=lxc
stop on net-down

emits net-device-up
emits static-network-up

console none

pre-start script
	if test -x "$(command -v ip)"; then
		ip addr add 127.0.0.1 dev lo || true
		ip -6 addr add ::1 dev lo || true
		ip link set lo up || true
	fi
        initctl emit -n net-device-up \
	    IFACE=lo LOGICAL=lo ADDRFAM=inet METHOD=loopback
	test -x "$(command -v ifup)" || \
	    initctl emit -n static-network-up
end script
