# net/lxcguestnet.conf
description "lxcguestnet: minimal networking in Linux Containers"

start on net-up and container CONTAINER=lxc
emits net-device-up
emits static-network-up

pre-start script
	ip addr add 127.0.0.1 dev lo || true
	ip -6 addr add ::1 dev lo || true
	ip link set lo up || true
        initctl emit -n net-device-up \
            IFACE=lo LOGICAL=lo ADDRFAM=inet METHOD=loopback || true
	initctl emit -n static-network-up
end script
