# net/neard.conf
description "neard - Near Field Communication (NFC) daemon"

start on started dbus
stop on stopping dbus

respawn

pre-start script
    test -x /usr/lib/neard/neard || { stop; exit 0; }
    modprobe -q nfc || true
end script

exec /usr/lib/neard/neard --nodaemon
