# net/tor.conf
description "tor - the Onion anonymity network daemon"

start on filesystem and static-network-up
stop on deconfiguring-networking or runlevel [016] or unmounting-filesystem

# TODO: readiness
respawn

# Stop tor from hogging the CPU with an appropiate nice level
#nice 5

env RUNDIR=/var/run/tor
env LOGDIR=/var/log/tor

pre-start script
    [ -x /usr/bin/tor ] || { stop; exit 0; }

    for dir in "$RUNDIR $LOGDIR"; do
        if [ ! -d $dir ]; then
            install -d -m 02750 -o debian-tor -g debian-tor $dir
            [ ! -x /sbin/restorecon ] || /sbin/restorecon $dir
        fi
    done

    /usr/bin/tor --verify-config || { stop; exit 0; }
end script

# TODO: nice values
# TODO: AA_EXEC stuff
# TODO: max fd stuff
exec /usr/bin/tor --hush
