# authentication/sssd.conf
description "System Security Services Daemon"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

respawn
exec /usr/sbin/sssd -f
