# authentication/sudo-ldap.conf
description "sudo - Provide limited superuser priveleges to specific users"

start on filesystem
stop on runlevel [016] or unmounting-filesystem

# Ensure that priveleges do not persist across reboots
pre-start script
	if test -d /var/lib/sudo; then
		find /var/lib/sudo -exec touch -d @0 '{}' \;
	fi
end script
